# sigma/rules/windows/process_creation/proc_creation_win_schtasks_creation.yml
title: Scheduled Task Creation Via Schtasks.EXE
id: 92626ddd-662c-49e3-ac59-f6535f12d189
status: test
description: Detects the creation of scheduled tasks by user accounts via the "schtasks"
  utility.
references:
- https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-create
author: Florian Roth (Nextron Systems)
date: 2019-01-16
modified: 2024-01-18
tags:
- attack.execution
- attack.persistence
- attack.privilege-escalation
- attack.t1053.005
- attack.s0111
- car.2013-08-001
- stp.1u
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|endswith: '\schtasks.exe'
    CommandLine|contains: ' /create '
  filter_main_system_user:
    User|contains:     # covers many language settings
    - 'AUTHORI'
    - 'AUTORI'
  condition: selection and not 1 of filter_main_*
falsepositives:
- Administrative activity
- Software installation
level: low
notes: |
  ### Technical Context
  The Sigma rule titled "Scheduled Task Creation Via Schtasks.EXE" is designed to detect when users create scheduled tasks using the Windows utility "schtasks.exe." This mechanism is often employed by both legitimate administrators and malicious actors to establish persistence on a compromised system. Specifically, the rule targets process creation events where the `Image` field indicates the use of "schtasks.exe" and the `CommandLine` includes the `/create` parameter. 
  The rule leverages process creation logs from Windows, allowing analysts to track actions taken by users that may indicate suspicious behavior, especially if executed by accounts that do not typically perform these actions. The associated MITRE ATT&CK technique is T1053.005, which relates to the creation of scheduled tasks for persistence, and the execution of code which could potentially lead to privilege escalation. This makes it a critical piece for understanding the potential establishment of lasting control by an adversary on a targeted system.
  ### Investigation Steps
  - **Review Process Creation Logs:** Examine process creation logs in your EDR tool to find the associated entries for `schtasks.exe`, focusing on unusual user accounts or command patterns that may indicate unauthorized scheduled task creation.
  - **Verify Scheduled Task Entries:** Use task scheduler or PowerShell (e.g., `Get-ScheduledTask`) to review the currently configured scheduled tasks and correlate them with recent process creation events.
  - **Investigate User Account Activity:** Analyze the activity and permissions of the user account that executed the `schtasks.exe` command; check for any irregularities in account behavior or elevated privileges.
  - **Check for Related Network Activity:** Query firewall and network logs for connections made by the system after the scheduled task creation event that may indicate further compromise or command-and-control communication.
  ### Prioritization
  This alert is considered low severity in an enterprise environment because administrative tasks, including the creation of scheduled tasks, are common and often legitimate. However, anomalous behavior around such activities merits further investigation to ensure no malicious intent is present.
  ### Blind Spots and Assumptions
  This rule may not fire if the scheduled task creation is executed via alternate methods or scripts that do not utilize `schtasks.exe`. Adversaries may employ other mechanisms to create persistence, such as PowerShell scripts or direct interactions with the registry, which would bypass this specific detection. Additionally, there is an underlying assumption that the telemetry from process creation is consistently collected and archived for adequate historical review, which may not be the case in all environments.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
