# sigma/rules/windows/process_creation/proc_creation_win_schtasks_creation_temp_folder.yml
title: Suspicious Scheduled Task Creation Involving Temp Folder
id: 39019a4e-317f-4ce3-ae63-309a8c6b53c5
status: test
description: Detects the creation of scheduled tasks that involves a temporary folder
  and runs only once
references:
- https://discuss.elastic.co/t/detection-and-response-for-hafnium-activity/266289/3
author: Florian Roth (Nextron Systems)
date: 2021-03-11
modified: 2022-10-09
tags:
- attack.execution
- attack.persistence
- attack.t1053.005
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|endswith: '\schtasks.exe'
    CommandLine|contains|all:
    - ' /create '
    - ' /sc once '
    - '\Temp\'
  condition: selection
falsepositives:
- Administrative activity
- Software installation
level: high
notes: |
  ### Technical Context
  This Sigma rule focuses on detecting the creation of suspicious scheduled tasks in Windows environments that involve the use of temporary folders and are configured to run only once. The rule specifically looks for instances where the `schtasks.exe` command-line tool is invoked with the `/create` and `/sc once` parameters, indicating that an administrator or a potentially malicious actor is creating a new scheduled task. The presence of the `\Temp\` folder in the command line raises concerns, as attackers often exploit temporary directories to execute malicious scripts in a stealthy manner. The pertinent threat techniques are captured under the MITRE ATT&CK framework, particularly the `T1053.005` technique for Scheduled Task/Job, which is relevant for both execution and persistence tactics.
  The technical data sources leveraged for this rule primarily include Windows process creation logs. These logs provide critical insights into when and how processes are initiated on the system, including the details of command-line arguments passed to the process. By analyzing these logs, incident responders can effectively identify and investigate potentially unauthorized or malicious task creation events.
  ### Investigation Steps
  - **Review Process Creation Logs**: Utilize EDR and Windows Event Logs to identify instances of `schtasks.exe` running, particularly focusing on entries related to the suspicious `/create`, `/sc once`, and `\Temp\` parameters.
  - **Analyze Scheduled Task Properties**: Use PowerShell or Task Scheduler on the affected system to examine the properties of the identified task(s) for abnormal execution times, associated commands, or indications of persistence mechanisms.
  - **Correlate with Other Data Sources**: Investigate any associated network traffic using NDR tools or endpoint indicators from firewalls to check for outbound connections initiated by tasks that may have been created, and check proxy logs for any relevant web requests.
  - **Check for Similar Events**: Look for patterns of similar behavior in recent Windows Event Logs for other task creations that might reflect coordinated or repeated malicious activity.
  ### Prioritization
  This alert is categorized as high severity because the creation of potentially malicious scheduled tasks can facilitate further persistence and execution of malware within an enterprise, posing significant risks to the integrity and confidentiality of systems and data.
  ### Blind Spots and Assumptions
  One potential blind spot of this rule is the reliance on the `Temp` directory naming convention, which might not universally apply. If attackers use different folder paths or if legitimate administrative tasks occur outside of this parameterization, those activities could go undetected. Additionally, the rule may not trigger in environments with more permissive administrative controls that regularly create scheduled tasks for maintenance purposes, leading to a higher potential for false positives. Also, if the attackers use obfuscation techniques in their command lines, it may defeat the rule.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality and importance, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
