# sigma/rules/windows/process_creation/proc_creation_win_reg_add_run_key.yml
title: Potential Persistence Attempt Via Run Keys Using Reg.EXE
id: de587dce-915e-4218-aac4-835ca6af6f70
status: test
description: Detects suspicious command line reg.exe tool adding key to RUN key in
  Registry
references:
- https://app.any.run/tasks/9c0f37bc-867a-4314-b685-e101566766d7/
- https://learn.microsoft.com/en-us/windows/win32/setupapi/run-and-runonce-registry-keys
author: Florian Roth (Nextron Systems)
date: 2021-06-28
modified: 2023-01-30
tags:
- attack.persistence
- attack.t1547.001
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    CommandLine|contains|all:
    - 'reg'
    - ' ADD '
    - 'Software\Microsoft\Windows\CurrentVersion\Run'
  condition: selection
falsepositives:
- Legitimate software automatically (mostly, during installation) sets up autorun
  keys for legitimate reasons.
- Legitimate administrator sets up autorun keys for legitimate reasons.
- Discord
level: medium
notes: |
  ### Technical Context
  This detection rule focuses on identifying potential persistence attempts by tracking the use of the Windows `reg.exe` utility, specifically in relation to adding entries to the RUN key in the Windows Registry. The RUN key, located at `Software\Microsoft\Windows\CurrentVersion\Run`, is a common method used by malicious actors to ensure that their payloads execute upon system startup. The rule captures process creation events where the command line contains specific phrases related to the use of `reg.exe` for adding a registry entry. The primary data sources utilized in this detection are process creation logs, which provide insight into command line arguments and the processes initiated across the Windows environment. This detection pertains to the MITRE ATT&CK framework under the tactic of Persistence and technique T1547.001.
  ### Investigation Steps
  - **Verify Process Origin:** Check the source of the `reg.exe` process using EDR tools to determine whether it originated from a legitimate location or an unusual directory.
  - **Examine Command Context:** Review the full command line arguments in the process creation logs to identify any other suspicious parameters or additional commands that may indicate malicious behavior.
  - **Query Registry Modifications:** Use Windows Registry auditing tools or commands to investigate the actual changes made to the RUN key and check for any other anomalous entries.
  - **Cross-Reference Alerts:** Correlate the findings with other security tools such as AV, Proxy, or Firewall logs to see if there were concurrent suspicious activities or attempts from the same source.
  ### Prioritization
  When the alert is fired, it indicates a medium severity level because while the creation of RUN keys is a common and sometimes legitimate action, it can also be indicative of unauthorized persistence mechanisms employed by malware or attackers.
  ### Blind Spots and Assumptions
  This rule may not fire if the command used to modify the RUN key does not leverage `reg.exe` or if it employs different methods to achieve persistence, such as using PowerShell scripts or other registry manipulation techniques. Additionally, legitimate administrative actions or software installations may generate false positives, leading to unnecessary investigations. Analysts should consider the context of the execution environment and historical patterns of registry modifications to differentiate between benign and malicious activities.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
