# sigma/rules/windows/process_creation/proc_creation_win_schtasks_schedule_via_masqueraded_xml_file.yml
title: Suspicious Scheduled Task Creation via Masqueraded XML File
id: dd2a821e-3b07-4d3b-a9ac-929fe4c6ca0c
status: test
description: Detects the creation of a scheduled task using the "-XML" flag with a
  file without the '.xml' extension. This behavior could be indicative of potential
  defense evasion attempt during persistence
references:
- https://learn.microsoft.com/en-us/windows/win32/taskschd/daily-trigger-example--xml-
- https://github.com/elastic/protections-artifacts/blob/084067123d3328a823b1c3fdde305b694275c794/behavior/rules/persistence_suspicious_scheduled_task_creation_via_masqueraded_xml_file.toml
author: Swachchhanda Shrawan Poudel, Elastic (idea)
date: 2023-04-20
modified: 2024-12-01
tags:
- attack.defense-evasion
- attack.persistence
- attack.t1036.005
- attack.t1053.005
logsource:
  product: windows
  category: process_creation
detection:
  selection_img:
  - Image|endswith: '\schtasks.exe'
  - OriginalFileName: 'schtasks.exe'
  selection_cli_create:
    CommandLine|contains:
    - '/create'
    - '-create'
  selection_cli_xml:
    CommandLine|contains:
    - '/xml'
    - '-xml'
  filter_main_extension_xml:
    CommandLine|contains: '.xml'
  filter_main_system_process:
    IntegrityLevel:
    - 'System'
    - 'S-1-16-16384'
  filter_main_rundll32:
    ParentImage|endswith: '\rundll32.exe'
    ParentCommandLine|contains|all:
    - ':\WINDOWS\Installer\MSI'
    - '.tmp,zzzzInvokeManagedCustomActionOutOfProc'
  filter_optional_third_party:
    ParentImage|endswith:
            # Consider removing any tools that you don't use to avoid blind spots
    - ':\ProgramData\OEM\UpgradeTool\CareCenter_*\BUnzip\Setup_msi.exe'
    - ':\Program Files\Axis Communications\AXIS Camera Station\SetupActions.exe'
    - ':\Program Files\Axis Communications\AXIS Device Manager\AdmSetupActions.exe'
    - ':\Program Files (x86)\Zemana\AntiMalware\AntiMalware.exe'
    - ':\Program Files\Dell\SupportAssist\pcdrcui.exe'
  condition: all of selection_* and not 1 of filter_main_* and not 1 of filter_optional_*
falsepositives:
- Unknown
level: medium
notes: |
  ### Technical Context
  This detection rule targets the creation of suspicious scheduled tasks on Windows systems, particularly those that utilize the `schtasks.exe` utility with the `-XML` flag, while using a file that lacks the standard '.xml' file extension. This behavior can imply an adversary attempting to evade defenses and establish persistence on a system, as it seeks to manipulate legitimate administrative tools for malicious purposes. The rule leverages telemetry from the process creation logs in Windows, specifically monitoring command line parameters to identify when scheduled tasks are created with command sequences that are consistent with this method of evasion. Moreover, it references relevant MITRE ATT&CK tactics and techniques such as `T1036.005` (Masquerading) and `T1053.005` (Scheduled Task/Job), highlighting the relevance of this behavior in the context of broader malware and persistence techniques.
  ### Investigation Steps
  - **Verify Scheduled Task Creation:** Check the Windows Task Scheduler for scheduled tasks created around the time of the alert, focusing on those that match the suspicious criteria noted in the alert and review task properties to identify any unusual configurations.
  - **Check Command Line Usage:** Utilize EDR to analyze the command line execution associated with `schtasks.exe`. Investigate if the command line includes unexpected options or parameters that suggest malicious intent.
  - **Examine Parent Processes:** Look into the parent processes of the `schtasks.exe` creation to determine the context in which the task was created, particularly focusing on instances where it was spawned by `rundll32.exe` or other known utilities that could indicate abnormal behavior.
  - **Cross-reference with Security Logs:** Review Windows Event Logs, especially Security and Sysmon logs, to identify other related activities occurring around the same time, such as unusual logins or privilege escalations that could provide context regarding the nature of the scheduled task.
  ### Prioritization
  The severity level of this alert is medium, as the creation of a scheduled task with deceptive file formats suggests a sophisticated attempt at persistence, which could lead to further malicious actions if not investigated promptly.
  ### Blind Spots and Assumptions
  This rule may not trigger if attackers utilize the expected '.xml' file extension or if they execute 'schtasks.exe' without using the `-XML` flag. Additionally, certain legitimate administrative or maintenance scripts that may also follow similar patterns might result in false positives. It is essential for detection engineers to consider these factors when tuning the detection logic, as well as the possibility of adversaries utilizing custom tools or scripts that are designed to evade such monitoring.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
