# sigma/rules/windows/process_creation/proc_creation_win_schtasks_guid_task_name.yml
title: Suspicious Scheduled Task Name As GUID
id: ff2fff64-4cd6-4a2b-ba7d-e28a30bbe66b
status: test
description: Detects creation of a scheduled task with a GUID like name
references:
- https://thedfirreport.com/2022/10/31/follina-exploit-leads-to-domain-compromise/
- https://thedfirreport.com/2022/02/21/qbot-and-zerologon-lead-to-full-domain-compromise/
author: Nasreddine Bencherchali (Nextron Systems)
date: 2022-10-31
tags:
- attack.execution
- attack.t1053.005
logsource:
  product: windows
  category: process_creation
detection:
  selection_img:
    Image|endswith: '\schtasks.exe'
    CommandLine|contains: '/Create '
  selection_tn:
    CommandLine|contains:
            # Can start with single or double quote
    - '/TN "{'
    - "/TN '{"
    - "/TN {"
  selection_end:
    CommandLine|contains:
            # Ending of the name to avoid possible FP in the rest of the commandline
    - '}"'
    - "}'"
    - '} '
  condition: all of selection_*
falsepositives:
- Legitimate software naming their tasks as GUIDs
level: medium
notes: |
  ### Technical Context
  This detection rule identifies the creation of suspicious scheduled tasks that have names formatted as GUIDs. In particular, it monitors for the use of the `schtasks.exe` command, which is commonly utilized to create, delete, configure, or display scheduled tasks in Windows environments. The rule focuses on the command line parameters that contain a task name (`/TN`) formatted as a GUID, specifically when the task name appears enclosed in braces (e.g., `{GUID}`). Given the potential for adversaries to create scheduled tasks to maintain persistence on compromised systems, this rule aligns with the MITRE ATT&CK technique T1053.005, which addresses the use of scheduled tasks for execution purposes. 
  The data sources leveraged in this rule include process creation logs from Windows, which capture details about the creation of processes, their command line arguments, and the executables used. By monitoring these specific indicators, the detection rule helps identify potentially malicious activities, offering an essential layer of threat detection within enterprise systems.
  ### Investigation Steps
  - **Verify Scheduled Task Creation:** Use your EDR solution to review the timeline of scheduled task creation events. Focus on recent entries that match the GUID naming pattern and assess their origins.
  - **Inspect Command Line Logging:** Examine the command line parameters captured in Windows event logs to ascertain the context in which `schtasks.exe` was executed. Determine if the execution context is legitimate or suspicious.
  - **Correlate with Other Alerts:** Look for related alerts in the SIEM from other telemetry sources that may indicate lateral movement or additional malicious activities associated with the identified task.
  - **Review Task Actions:** Investigate what actions the created task is configured to perform. Assess whether these actions align with normal behavior or exhibit malicious intentions (e.g., downloading files, executing scripts).
  ### Prioritization
  The medium severity level is assigned due to the potential for threat actor presence within the environment. The creation of scheduled tasks with GUIDs can signify attempts to maintain persistence and necessitates prompt investigation to determine the extent of any compromises.
  ### Blind Spots and Assumptions
  This rule may not fire in environments where legitimate software or administrative personnel create scheduled tasks with GUID-like names, potentially leading to false positives. Additionally, if the task name doesn't conform to the expected GUID format or is obfuscated, the rule might miss such activities. Adversaries could also deploy alternate techniques to launch persistence mechanisms that do not use the conventional `schtasks.exe` command or follow similar patterns. Awareness of these factors and continuous refinement of detection logic is essential for effective monitoring.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
