# sigma/rules/windows/process_creation/proc_creation_win_schtasks_system.yml
title: Schtasks Creation Or Modification With SYSTEM Privileges
id: 89ca78fd-b37c-4310-b3d3-81a023f83936
status: test
description: Detects the creation or update of a scheduled task to run with "NT AUTHORITY\SYSTEM"
  privileges
references:
- https://www.elastic.co/security-labs/exploring-the-qbot-attack-pattern
- https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks
author: Nasreddine Bencherchali (Nextron Systems)
date: 2022-07-28
modified: 2025-02-15
tags:
- attack.execution
- attack.persistence
- attack.t1053.005
logsource:
  product: windows
  category: process_creation
detection:
  selection_root:
    Image|endswith: '\schtasks.exe'
    CommandLine|contains:
    - ' /change '
    - ' /create '
  selection_run:
    CommandLine|contains: '/ru '
  selection_user:
    CommandLine|contains:
    - 'NT AUT'         # This covers the usual NT AUTHORITY\SYSTEM
    - ' SYSTEM '         # SYSTEM is a valid value for schtasks hence it gets it's own value with space
  filter_optional_teamviewer:
        # FP from test set in SIGMA
        # Cannot use ParentImage on all OSes for 4688 events
        # ParentImage|contains|all:
        #     - '\AppData\Local\Temp\'
        #     - 'TeamViewer_.exe'
    Image|endswith: '\schtasks.exe'
    CommandLine|contains|all:
    - '/TN TVInstallRestore'
    - '\TeamViewer_.exe'
  filter_optional_office:
    CommandLine|contains|all:
            # https://answers.microsoft.com/en-us/msoffice/forum/all/office-15-subscription-heartbeat-task-created-on/43ab5e53-a9fb-47c6-8c14-44889974b9ff
    - 'Subscription Heartbeat'
    - '\HeartbeatConfig.xml'
    - '\Microsoft Shared\OFFICE'
  filter_optional_avira:
    CommandLine|contains:
    - '/Create /F /RU System /SC WEEKLY /TN AviraSystemSpeedupVerify /TR '
    - ':\Program Files (x86)\Avira\System Speedup\setup\avira_speedup_setup.exe'
    - '/VERIFY /VERYSILENT /NOSTART /NODOTNET /NORESTART" /RL HIGHEST'
  condition: all of selection_* and not 1 of filter_optional_*
falsepositives:
- Unknown
level: high
notes: |
  ### Technical Context
  The Sigma rule titled "Schtasks Creation Or Modification With SYSTEM Privileges" is designed to detect the creation or modification of scheduled tasks that execute with elevated privileges—specifically, those associated with the "NT AUTHORITY\SYSTEM" account. This type of activity is commonly associated with adversaries attempting to establish persistence on a system by leveraging powerful scheduling capabilities inherent to Windows. The rule monitors the command-line inputs of the `schtasks.exe` process for certain key indicators such as the commands `/create` or `/change`, coupled with the use of the `NT AUTHORITY\SYSTEM` user context. 
  The relevant data sources include process creation logs that track the execution of `schtasks.exe`, which is logged as part of Windows event logging. This aligns with the MITRE ATT&CK technique T1053.005 (Scheduled Task/Job: Scheduled Task Creation), which details how malicious actors may create tasks as a method of maintaining persistence.
  ### Investigation Steps
  - **Review Process Creation Logs:** Utilize your EDR tool to examine the process creation logs specifically for `schtasks.exe` executions, paying attention to the command-line parameters used during the task creation or modification.
  - **Correlate with User Privileges:** Check the associated account privileges in the Windows Security logs to confirm whether the tasks were indeed created under the "NT AUTHORITY\SYSTEM" context.
  - **Monitor for Anomalous Tasks:** Investigate any newly created or modified scheduled tasks that may not align with organizational policy, particularly those originating from unusual or unauthorized user accounts.
  - **Examine Relevant Alerts:** Search in the SIEM for alerts related to process creation and task scheduling that coincide with known malicious activity patterns, especially if TeamViewer or Office tasks are involved as potential benign filters.
  ### Prioritization
  This alert is classified as high severity due to the implications of a potentially malicious actor gaining SYSTEM-level privileges, which can enable extensive control over the affected system and facilitate further exploitation or lateral movement.
  ### Blind spots and Assumptions
  One limitation of this rule is the reliance on process creation logs, which may not capture all instances of legitimate scheduled task modifications or creations, especially in environments with heavy use of automation tools. Additionally, adversaries may use alternative methods to create tasks, such as through PowerShell or directly manipulating Windows APIs, which would not trigger this rule. The potential for false positives exists if legitimate applications (like TeamViewer) perform similar operations under the SYSTEM account. Therefore, continuous tuning and contextual analysis of alerts are recommended.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
