# sigma/rules/windows/process_creation/proc_creation_win_cmd_assoc_execution.yml
title: Change Default File Association Via Assoc
id: 3d3aa6cd-6272-44d6-8afc-7e88dfef7061
related:
- id: ae6f14e6-14de-45b0-9f44-c0986f50dc89
  type: similar
status: test
description: |
  Detects file association changes using the builtin "assoc" command.
  When a file is opened, the default program used to open the file (also called the file association or handler) is checked. File association selections are stored in the Windows Registry and can be edited by users, administrators, or programs that have Registry access or by administrators using the built-in assoc utility. Applications can modify the file association for a given file extension to call an arbitrary program when a file with the given extension is opened.
references:
- https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1546.001/T1546.001.md
author: Timur Zinniatullin, oscd.community
date: 2019-10-21
modified: 2023-03-06
tags:
- attack.persistence
- attack.t1546.001
logsource:
  category: process_creation
  product: windows
detection:
  selection_img:
  - Image|endswith: '\cmd.exe'
  - OriginalFileName: 'Cmd.Exe'
  selection_cli:
    CommandLine|contains: 'assoc'
  condition: all of selection_*
fields:
- Image
- CommandLine
- User
- LogonGuid
- Hashes
- ParentProcessGuid
- ParentCommandLine
falsepositives:
- Admin activity
level: low
notes: |
  n
  ### Technical Context
  This detection rule monitors for changes to file associations in the Windows operating system using the built-in "assoc" command. File associations determine which application is used to open specific file types, and these settings are stored in the Windows Registry. By modifying file associations, an attacker could potentially associate malicious applications with file types, thus enabling their execution without user consent. The rule captures process creation events where the command executed includes "assoc," specifically looking for instances where the command is invoked through `cmd.exe`. The relevant technical data sources underpinning this rule are Windows Process Creation logs and the Sysmon utility, typically focusing on attributes like the image name, command line contents, and user context. This detection aligns with the MITRE ATT&CK technique T1546.001 (Event Triggered Execution: File Association Modification) under the persistence tactic.
  ### Investigation Steps
  - **Check Process Context:** Review the process creation logs in your EDR to identify the user and logon context under which the "assoc" command was executed. Correlate this with user activity to assess legitimacy.
  - **Review Parent Processes:** Identify parent processes associated with the invocation of the "assoc" command to determine if it was launched by a trusted application or if there's a possibility of process tampering.
  - **Audit Registry Changes:** Use Registry monitoring tools to identify any recent changes to file associations in the Windows Registry. Look specifically for modifications related to the file types that were altered.
  - **Examine User Behavior:** Investigate the user's history surrounding this event, including recent actions and whether they could have legitimately modified file associations as part of their duties.
  ### Prioritization
  The alert is assigned a low severity level; however, it raises concern due to the potential for misuse in modifying file associations which can lead to malicious software execution, particularly if the command was run under suspicious circumstances or by an unauthorized user.
  ### Blind Spots and Assumptions
  This rule may not fire in environments where administrative actions are routine and expected, possibly resulting in legitimate use cases from authorized personnel that could dilute the signal of malicious activity. Additionally, if an attacker gains persistence through other means (e.g., scheduled tasks), they may bypass this detection entirely. Furthermore, the execution of the "assoc" command might not be logged if Sysmon is not appropriately configured or if Process Creation logging is disabled.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
