# sigma/rules/windows/process_creation/proc_creation_win_schtasks_one_time_only_midnight_task.yml
title: Uncommon One Time Only Scheduled Task At 00:00
id: 970823b7-273b-460a-8afc-3a6811998529
status: test
description: Detects scheduled task creation events that include suspicious actions,
  and is run once at 00:00
references:
- https://www.trendmicro.com/vinfo/us/security/news/ransomware-spotlight/ransomware-spotlight-blackbyte
author: pH-T (Nextron Systems)
date: 2022-07-15
modified: 2023-02-03
tags:
- attack.execution
- attack.persistence
- attack.privilege-escalation
- attack.t1053.005
logsource:
  category: process_creation
  product: windows
detection:
  selection_img:
  - Image|contains: '\schtasks.exe'
  - OriginalFileName: 'schtasks.exe'
  selection_cli:
    CommandLine|contains:
    - 'wscript'
    - 'vbscript'
    - 'cscript'
    - 'wmic '
    - 'wmic.exe'
    - 'regsvr32.exe'
    - 'powershell'
    - '\AppData\'
  selection_time:
    CommandLine|contains|all:
    - 'once'
    - '00:00'
  condition: all of selection_*
falsepositives:
- Software installation
level: high
notes: |
  n
  ### Technical Context
  This detection rule identifies the creation of scheduled tasks that execute suspicious commands at the specific time of 00:00 (midnight). It specifically looks for the use of the `schtasks.exe` utility, commonly employed by attackers to create persistent and stealthy task execution mechanisms. The rule targets command line parameters that include various scripting engines, as well as common tools that might be leveraged by adversaries, such as Windows Management Instrumentation Command-line (WMIC), `regsvr32.exe`, and PowerShell. The MITRE ATT&CK framework tags associated with this rule fall under tactics related to **Execution** (T1053.005 for creating scheduled tasks) and **Persistence**, indicating techniques used by adversaries for maintaining access to a system.
  The telemetry data involved in this detection primarily comes from Windows Process Creation logs, which include command-line arguments and executable paths, allowing the detection of these potentially malicious scheduled task creation events. By focusing on both the command line and the specific time frame for execution, the rule aims to reduce false positives from legitimate scheduled tasks, thus enhancing the security monitoring capabilities against sophisticated threats that may attempt to blend in with normal operational procedures.
  ### Investigation Steps
  - **Review Scheduled Task History:** Leverage the Task Scheduler on the affected endpoint to investigate the details of the created scheduled tasks, focusing on the execution time and the commands configured to run. 
  - **Analyze Command-Line Context:** Utilize EDR tools to review the command-line history associated with `schtasks.exe` executions to confirm the legitimacy of the arguments used and their purpose.
  - **Identify Model Activities:** Review Windows Event Logs for other events around the same timestamp to understand if there were any preceding or following suspicious activities, such as unusual logins or process creations.
  - **Check for User Context:** Determine the user account under which the scheduled task was created and verify if this account has administrative privileges that may facilitate persistence or privilege escalation.
  ### Prioritization
  The alert generated by this rule is classified as high severity due to the potential implications of a scheduled task being created to run malicious scripts or commands, particularly at a time when system activity is likely low. This behavior may indicate an attempt at stealth or persistence within the environment, warranting immediate investigation.
  ### Blind Spots and Assumptions
  This detection rule may produce false negatives if an attacker uses legitimate tasks or familiar scheduling names outside the specified criteria, or if script execution is accomplished through less common methods. Additionally, if an adversary were to mask their actions using legitimate user accounts or during atypical times or utilize tools that do not invoke `schtasks.exe`, the rule may not fire. Furthermore, environments with highly customized task configurations may also evade detection if they don’t match the expected parameters outlined in this rule.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
