# sigma/rules/windows/process_creation/proc_creation_win_schtasks_persistence_windows_telemetry.yml
title: Potential Persistence Via Microsoft Compatibility Appraiser
id: f548a603-c9f2-4c89-b511-b089f7e94549
related:
- id: 73a883d0-0348-4be4-a8d8-51031c2564f8
  type: derived
status: test
description: |
  Detects manual execution of the "Microsoft Compatibility Appraiser" task via schtasks.
  In order to trigger persistence stored in the "\AppCompatFlags\TelemetryController" registry key.
references:
- https://www.trustedsec.com/blog/abusing-windows-telemetry-for-persistence/
author: Sreeman
date: 2020-09-29
modified: 2023-02-10
tags:
- attack.persistence
- attack.t1053.005
logsource:
  product: windows
  category: process_creation
detection:
  selection_img:
  - Image|endswith: '\schtasks.exe'
  - OriginalFileName: 'schtasks.exe'
  selection_cli:
    CommandLine|contains|all:
    - 'run '
    - '\Application Experience\Microsoft Compatibility Appraiser'
  condition: all of selection_*
falsepositives:
- Unknown
level: medium
notes: |
  ### Technical Context
  This Sigma rule focuses on detecting the potentially malicious use of the Microsoft Compatibility Appraiser as a mechanism for persistence within a Windows environment. The rule specifically looks for processes related to the execution of the schtasks.exe utility, which is a legitimate Windows command-line tool used for scheduling tasks. This involves monitoring the process creation logs to identify instances where schtasks.exe is invoked with a command line that includes the "Microsoft Compatibility Appraiser" task. This can indicate an attacker’s attempt to exploit a legitimate Windows feature for malicious purposes, by creating persistence in the system through registry modifications in the "\AppCompatFlags\TelemetryController" registry key. The rule aligns with the MITRE ATT&CK framework under the Tactic of Persistence, specifically the Technique T1053.005, which relates to Scheduled Task/Job.
  ### Investigation Steps
  - **Review EDR Alerts:** Check EDR for any suspicious execution of schtasks.exe, particularly focusing on the command line parameters that involve the Microsoft Compatibility Appraiser.
  - **Examine Process Creation Logs:** Look into process creation logs in the Windows Security logs to gather information on the parent process and user context for the execution of schtasks.exe.
  - **Investigate Registry Changes:** Query the registry for recent modifications to the "\AppCompatFlags\TelemetryController" key, verifying if any changes reflect unauthorized or suspicious activity.
  - **Analyze Network Connections:** Inspect network logs for any outbound connections made immediately after the identified task execution, which could indicate additional malicious activity or data exfiltration.
  ### Prioritization
  This alert is prioritized as medium severity due to the elevated risk of persistence enabled by exploiting benign Windows features, which could lead to prolonged access by adversaries and complicate incident response efforts.
  ### Blind Spots and Assumptions
  This rule may not fire in scenarios where an adversary has modified the schtasks.exe binary or if the command is executed with alternate command-line utilities that do not invoke schtasks directly. Additionally, the rule assumes that logging configurations are set to capture relevant process creation events and that the integrity of the telemetry registry key is intact.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
