# sigma/rules/windows/process_creation/proc_creation_win_schtasks_change.yml
title: Suspicious Modification Of Scheduled Tasks
id: 1c0e41cd-21bb-4433-9acc-4a2cd6367b9b
related:
- id: 614cf376-6651-47c4-9dcc-6b9527f749f4     # Security-Audting Eventlog
  type: similar
status: test
description: |
  Detects when an attacker tries to modify an already existing scheduled tasks to run from a suspicious location
  Attackers can create a simple looking task in order to avoid detection on creation as it's often the most focused on
  Instead they modify the task after creation to include their malicious payload
references:
- Internal Research
- https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks
author: Nasreddine Bencherchali (Nextron Systems)
date: 2022-07-28
modified: 2022-11-18
tags:
- attack.execution
- attack.t1053.005
logsource:
  product: windows
  category: process_creation
detection:
  selection_schtasks:
    Image|endswith: '\schtasks.exe'
    CommandLine|contains|all:
    - ' /Change '
    - ' /TN '
  selection_susp_locations:
    CommandLine|contains:
    - '\AppData\Local\Temp'
    - '\AppData\Roaming\'
    - '\Users\Public\'
    - '\WINDOWS\Temp\'
    - '\Desktop\'
    - '\Downloads\'
    - '\Temporary Internet'
    - 'C:\ProgramData\'
    - 'C:\Perflogs\'
    - '%ProgramData%'
    - '%appdata%'
    - '%comspec%'
    - '%localappdata%'
  selection_susp_images:
    CommandLine|contains:
    - 'regsvr32'
    - 'rundll32'
    - 'cmd /c '
    - 'cmd /k '
    - 'cmd /r '
    - 'cmd.exe /c '
    - 'cmd.exe /k '
    - 'cmd.exe /r '
    - 'powershell'
    - 'mshta'
    - 'wscript'
    - 'cscript'
    - 'certutil'
    - 'bitsadmin'
    - 'bash.exe'
    - 'bash '
    - 'scrcons'
    - 'wmic '
    - 'wmic.exe'
    - 'forfiles'
    - 'scriptrunner'
    - 'hh.exe'
    - 'hh '
  condition: all of selection_*
falsepositives:
- Unknown
level: high
notes: |
  n
  ### Technical Context
  The "Suspicious Modification of Scheduled Tasks" Sigma rule is designed to identify potential malicious activities involving modifications to scheduled tasks on Windows systems. Attackers often manipulate existing tasks instead of creating new ones to evade detection, using the `schtasks.exe` utility. The rule specifically looks for command-line arguments indicating a change (`/Change`) to tasks, while also monitoring for modifications made from suspicious locations such as user directories and system temporary folders. The command-line inputs with known executable patterns like `powershell`, `cmd.exe`, or `wmic` are further examined to gauge the legitimacy of the changes. This is relevant to the MITRE ATT&CK tactic of Execution, particularly technique T1053.005 (Scheduled Task/Job Modification).
  ### Investigation Steps
  - **Examine EDR Alerts**: Review EDR alerts to identify the specific process that triggered this rule, noting the user context and time of the modification.
  - **Analyze Command-Line Details**: Inspect command-line parameters used during the scheduled task modification to identify any suspicious executables or parameters.
  - **Review Scheduled Tasks**: Use the `schtasks` command to list current scheduled tasks and verify the integrity of tasks modified around the identified timeframe.
  - **Correlate with User Activity**: Check the Windows Security Event Logs for related user activity to determine if the user account involved in the modification has a history of unusual behavior.
  ### Prioritization
  Given that this alert indicates the potential manipulation of scheduled tasks, which can be leveraged to execute malware persistently and evade detection, it is considered a high-severity alert in an enterprise environment.
  ### Blind Spots and Assumptions
  There are potential blind spots where legitimate administrative activities might trigger this rule (e.g., scheduled maintenance tasks by system administrators). Furthermore, advanced adversaries may use encryption or obfuscation techniques in their command-line parameters to avoid detection. The assumption is that not all legitimate uses of `schtasks.exe` are malicious, requiring context for accurate assessment. 
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
