# sigma/rules/windows/process_creation/proc_creation_win_schtasks_schedule_type.yml
title: Suspicious Schtasks Schedule Types
id: 24c8392b-aa3c-46b7-a545-43f71657fe98
related:
- id: 7a02e22e-b885-4404-b38b-1ddc7e65258a
  type: similar
status: test
description: Detects scheduled task creations or modification on a suspicious schedule
  type
references:
- https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-change
- https://learn.microsoft.com/en-us/windows-server/administration/windows-commands/schtasks-create
- http://blog.talosintelligence.com/2022/09/lazarus-three-rats.html
author: Nasreddine Bencherchali (Nextron Systems)
date: 2022-09-09
tags:
- attack.execution
- attack.t1053.005
logsource:
  product: windows
  category: process_creation
detection:
  selection_img:
  - Image|endswith: '\schtasks.exe'
  - OriginalFileName: 'schtasks.exe'
  selection_time:
    CommandLine|contains:
    - ' ONLOGON '
    - ' ONSTART '
    - ' ONCE '
    - ' ONIDLE '
  filter_privs:
    CommandLine|contains:
    - 'NT AUT'         # This covers the usual NT AUTHORITY\SYSTEM
    - ' SYSTEM'         # SYSTEM is a valid value for schtasks hence it gets it's own value with space
    - 'HIGHEST'
  condition: all of selection_* and not 1 of filter_*
falsepositives:
- Legitimate processes that run at logon. Filter according to your environment
level: high
notes: |
  ### Technical Context
  This Sigma rule aims to detect the creation or modification of scheduled tasks in Windows environments, specifically identifying those that operate on unusual schedules that may signal malicious intent. The rule monitors process creation events related to the `schtasks.exe` utility, which is used to create, delete, configure, or display scheduled tasks on a Windows system. It specifically looks for command-line parameters that point to schedules such as `ONLOGON`, `ONSTART`, `ONCE`, and `ONIDLE`, which are less common and could be leveraged by adversaries for persistence or malicious activity. This can be particularly dangerous when combined with elevated privileges or unexpected users being used to run the task. The rule is aligned with the MITRE ATT&CK framework under the Execution tactic, specifically technique T1053.005, which deals with scheduled tasks.
  ### Investigation Steps
  - **Check Process Creation Logs:** Utilize EDR tools to search for `schtasks.exe` process creation events within the relevant time frame to gather contextual information about the task's creation or modification.
  - **Examine Command-Line Parameters:** Analyze the command-line arguments associated with the suspicious `schtasks` entries, particularly focusing on the presence of unusual triggers like `ONLOGON`, `ONSTART`, or other atypical options.
  - **Validate User Context:** Cross-reference the user context under which the scheduled task was created or modified, paying close attention to any non-administrative accounts or unexpected user accounts involved.
  - **Review Task Details & Execution Logs:** Use Windows Task Scheduler to review the details of the created tasks, confirm their intended purpose, and check for any corresponding execution logs for anomalies or unexpected behavior.
  ### Prioritization
  The alert generated by this rule is deemed high severity as it indicates the possibility of malicious persistence mechanisms being established on the system. Given the potential for exploitation, a prompt investigation is essential to mitigate any risks.
  ### Blind Spots and Assumptions
  This rule may not fire in situations where legitimate administrative operations are performed or if the tasks are created through alternate, less detectable means (e.g., scripts or alternative tools). Additionally, if an adversary uses obfuscation or disguises their activities by modifying command-line arguments used with `schtasks.exe`, it could evade detection. Analysts should remain cautious of potential false positives arising from legitimate environments where automated tasks are a normal operation.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
