# sigma/rules/windows/process_creation/proc_creation_win_reg_direct_asep_registry_keys_modification.yml
title: Direct Autorun Keys Modification
id: 24357373-078f-44ed-9ac4-6d334a668a11
status: test
description: Detects direct modification of autostart extensibility point (ASEP) in
  registry using reg.exe.
references:
- https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1547.001/T1547.001.md
author: Victor Sergeev, Daniil Yugoslavskiy, oscd.community
date: 2019-10-25
modified: 2022-08-04
tags:
- attack.persistence
- attack.t1547.001
logsource:
  category: process_creation
  product: windows
detection:
  selection_1:
    Image|endswith: '\reg.exe'
    CommandLine|contains: 'add'         # to avoid intersection with discovery tactic rules
  selection_2:
    CommandLine|contains:               # need to improve this list, there are plenty of ASEP reg keys
    - '\software\Microsoft\Windows\CurrentVersion\Run'         # Also covers the strings "RunOnce", "RunOnceEx", "RunServices", "RunServicesOnce"
    - '\software\Microsoft\Windows NT\CurrentVersion\Winlogon\Userinit'
    - '\software\Microsoft\Windows NT\CurrentVersion\Winlogon\Shell'
    - '\software\Microsoft\Windows NT\CurrentVersion\Windows'
    - '\software\Microsoft\Windows\CurrentVersion\Explorer\User Shell Folders'
    - '\system\CurrentControlSet\Control\SafeBoot\AlternateShell'
  condition: all of selection_*
fields:
- CommandLine
- ParentCommandLine
falsepositives:
- Legitimate software automatically (mostly, during installation) sets up autorun
  keys for legitimate reasons.
- Legitimate administrator sets up autorun keys for legitimate reasons.
- Discord
level: medium
notes: |
  ### Technical Context
  The "Direct Autorun Keys Modification" Sigma rule aims to detect unauthorized changes to critical autostart extensibility points (ASEP) in the Windows registry using the `reg.exe` command-line tool. This modification typically occurs when an attacker seeks to establish persistence on a compromised system by manipulating registry keys that automatically execute applications upon user login or system startup. Specifically, the rule monitors for process creation events related to `reg.exe` with command-line parameters indicating an addition of registry keys that fall under notable paths like `Run`, `Userinit`, and others associated with automatic execution, thereby indicating potential malicious activities aligned with the MITRE ATT&CK technique T1547.001 (Boot or Logon Initialization Scripts). The sources of telemetry involved in this detection primarily include Windows process creation logs alongside specific command-line parameters to pinpoint malicious behavior.
  ### Investigation Steps
  - **Check Process Creation Logs**: Use your EDR to review process creation logs around the time of the alert for any `reg.exe` instances, paying particular attention to their parent processes for context.
  - **Analyze Command-Line Arguments**: Examine the command-line details captured in the alert to identify any suspicious or unexpected registry keys being modified, especially those related to autorun settings.
  - **Validate System Configuration**: Query the targeted registry keys directly through your system management tools to confirm if the changes observed in the alerts were authorized or part of a legitimate application installation.
  - **Review User Activity Logs**: Inspect user activity around the time of the modification, such as local logins or administrative actions, to understand if legitimate administrative actions coincided with the detection.
  ### Prioritization
  This alert is classified as medium severity because unauthorized modifications to autorun keys can indicate persistent threats within the environment, potentially allowing an adversary to maintain access even after initial compromise attempts.
  ### Blind spots and Assumptions
  This rule may not fire in scenarios where an attacker uses scripts or applications that do not interact with `reg.exe` to alter autorun keys. Additionally, legitimate software installations or system configurations performed by administrators could trigger false positives, as they may also lead to changes in these registry keys. Furthermore, if the monitoring or logging for process creation events is disabled or misconfigured, the rule will not be able to capture relevant data.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
