# sigma/rules/windows/process_creation/proc_creation_win_schtasks_folder_combos.yml
title: Schtasks From Suspicious Folders
id: 8a8379b8-780b-4dbf-b1e9-31c8d112fefb
status: test
description: Detects scheduled task creations that have suspicious action command
  and folder combinations
references:
- https://symantec-enterprise-blogs.security.com/blogs/threat-intelligence/lazarus-dream-job-chemical
author: Florian Roth (Nextron Systems)
date: 2022-04-15
modified: 2022-11-18
tags:
- attack.execution
- attack.t1053.005
logsource:
  product: windows
  category: process_creation
detection:
  selection_img:
  - Image|endswith: '\schtasks.exe'
  - OriginalFileName: 'schtasks.exe'
  selection_create:
    CommandLine|contains: ' /create '
  selection_command:
    CommandLine|contains:
    - 'powershell'
    - 'pwsh'
    - 'cmd /c '
    - 'cmd /k '
    - 'cmd /r '
    - 'cmd.exe /c '
    - 'cmd.exe /k '
    - 'cmd.exe /r '
  selection_all_folders:
    CommandLine|contains:
    - 'C:\ProgramData\'
    - '%ProgramData%'
  condition: all of selection_*
falsepositives:
- Unknown
level: high
notes: |
  ### Technical Context
  The rule "Schtasks From Suspicious Folders" aims to detect the creation of scheduled tasks in Windows environments, particularly when these tasks are associated with potentially malicious command actions and suspicious folder paths. The primary focus is on the `schtasks.exe` executable, which is used to create, delete, configure, or display scheduled tasks. The rule triggers when the command line for creating the scheduled task includes suspicious strings indicative of obfuscation or abuse of legitimate functionality, such as invoking PowerShell or CMD commands. 
  Key telemetry for this detection stems from process creation logs, where it actively scrutinizes the command-line arguments for the creation of tasks that leverage paths like `C:\ProgramData\`, a common location for malicious actors to hide their activities. This detection relates to the MITRE ATT&CK framework under the Execution tactic, specifically technique T1053.005, which addresses scheduled task/cron job abuse.
  ### Investigation Steps
  - **Verify Scheduled Tasks:** Check the Windows Task Scheduler for newly created tasks that correspond to the alert, paying close attention to their configuration and whether they point to legitimate processes.
  - **Analyze Command Line Inputs:** Review the command line associated with the detected `schtasks.exe` creation to identify any suspicious arguments or paths that may indicate malicious intent.
  - **Cross-reference with EDR Data:** Utilize EDR tools to inspect the behavior of processes related to the scheduled tasks, looking for additional indicators of compromise such as unexpected executions or network calls.
  - **Check Event Logs for History:** Review Windows Event Logs (particularly Security, System, and Application) for any prior related events before the task creation, as this may provide context or a timeline leading to the alert.
  ### Prioritization
  This alert is classified as high severity due to the potential use of scheduled tasks for persistence or to execute malicious scripts, both of which can lead to wider compromise if not addressed promptly.
  ### Blind Spots and Assumptions
  This rule may not trigger if an attacker employs unique command-line syntax or renames legitimate command-line utilities to evade detection. It assumes that `schtasks.exe` is used for task creation, while other methods (like direct Registry manipulation) might be overlooked. Additionally, environments with heavy customizations or non-standard folder structures might hide malicious activities more effectively than usual.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
