# sigma/rules/windows/process_creation/proc_creation_win_at_interactive_execution.yml
title: Interactive AT Job
id: 60fc936d-2eb0-4543-8a13-911c750a1dfc
status: test
description: Detects an interactive AT job, which may be used as a form of privilege
  escalation.
references:
- https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1053.002/T1053.002.md
- https://eqllib.readthedocs.io/en/latest/analytics/d8db43cf-ed52-4f5c-9fb3-c9a4b95a0b56.html
author: E.M. Anhaus (originally from Atomic Blue Detections, Endgame), oscd.community
date: 2019-10-24
modified: 2021-11-27
tags:
- attack.privilege-escalation
- attack.t1053.002
logsource:
  category: process_creation
  product: windows
detection:
  selection:
    Image|endswith: '\at.exe'
    CommandLine|contains: 'interactive'
  condition: selection
falsepositives:
- Unlikely (at.exe deprecated as of Windows 8)
level: high
notes: |
  ### Technical Context
  The "Interactive AT Job" Sigma rule is designed to identify the execution of the AT (Scheduled Tasks) command-line utility with the "interactive" flag. This detection focus is pertinent because the AT utility can be exploited by malicious actors to execute commands with elevated privileges or to schedule tasks that may grant them unauthorized access, thereby allowing for privilege escalation. The primary data sources for this detection include process creation logs which capture relevant details about running processes, including the image name (specifically, `at.exe`) and the associated command-line parameters that indicate whether the job is intended to run interactively. This rule aligns with the MITRE ATT&CK framework under the privilege escalation tactic, specifically technique T1053.002, which deals with task scheduling methods used by attackers to escalate their privileges.
  ### Investigation Steps
  - **Check EDR for Context:** Review endpoint detection and response (EDR) alerts for any corresponding instances of `at.exe` being executed. Analyze the process tree to understand any parent-child relationships that led to this process being launched.
  - **Review Command-Line Parameters:** Examine the command-line parameters associated with the process execution to confirm if "interactive" was specified. Look for any unusual commands or unexpected scheduling tasks that could indicate malicious intent.
  - **Inspect Windows Security Logs:** Investigate Windows Security logs for Event ID 4698, which details the creation of a scheduled task. This can provide insights into when and by whom the task was created and whether it follows normal operational patterns.
  - **Network Activity Analysis:** Utilize network detection and response (NDR) tools to identify any associated network connections that may have occurred after the execution of the interactive task. Analyze if there were any outbound connections indicative of further malicious actions or command and control (C2) communications.
  ### Prioritization
  Given that the alert pertains to privileged task scheduling activity, it is classified with a high severity level. This classification is warranted because it can indicate an active attempt by an adversary to gain higher access on potentially sensitive systems, necessitating immediate investigation to mitigate potential threats.
  ### Blind Spots and Assumptions
  This detection rule may not trigger in environments where the `at.exe` tool has been removed or replaced, especially since `at.exe` has been deprecated since Windows 8, shifting to Task Scheduler for task management. Additionally, if an adversary uses alternate methods for privilege escalation or task scheduling that do not leverage `at.exe`, such as PowerShell or Windows Task Scheduler with renamed scripts, this rule may not capture these events. Therefore, it's critical for detection engineers to consider integrating additional telemetry and establish alerting for a wider string of potential privilege escalation techniques.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
