# sigma/rules/windows/registry/registry_event/registry_event_runkey_winekey.yml
title: WINEKEY Registry Modification
id: b98968aa-dbc0-4a9c-ac35-108363cbf8d5
status: test
description: Detects potential malicious modification of run keys by winekey or team9
  backdoor
references:
- https://www.fireeye.com/blog/threat-research/2020/10/kegtap-and-singlemalt-with-a-ransomware-chaser.html
author: omkar72
date: 2020-10-30
modified: 2021-11-27
tags:
- attack.persistence
- attack.t1547
logsource:
  category: registry_event
  product: windows
detection:
  selection:
    TargetObject|endswith: 'Software\Microsoft\Windows\CurrentVersion\Run\Backup Mgr'
  condition: selection
fields:
- ComputerName
- Image
- EventType
- TargetObject
falsepositives:
- Unknown
level: high
notes: |
  ### Technical Context
  The "WINEKEY Registry Modification" Sigma rule is designed to detect unauthorized modifications to Windows registry run keys, specifically targeting the `Backup Mgr` entry under `HKEY_CURRENT_USER\Software\Microsoft\Windows\CurrentVersion\Run`. This registry location is commonly used by malicious actors to achieve persistence on a Windows system by ensuring that malware is executed at user logon. The rule analyzes registry event logs, specifically looking for changes made to the `TargetObject` that match the specified path, and is particularly focused on changes that may originate from tools associated with known threat actors, such as Team9.
  This detection aligns with the MITRE ATT&CK tactic of **Persistence** (T1547), indicating that it aims to uncover methods adversaries employ to maintain a foothold within a compromised environment. By monitoring registry events, the rule leverages Windows Security Logs to detect possible malicious activities that may not be evident through standard application behavior.
  ### Investigation Steps
  - **Check Initial Alerts and Context**: Investigate the triggered alert in the SIEM, reviewing the event's timestamp, computer name, and user account involved for context.
  - **Correlate with EDR Data**: Use your EDR to review process activity around the time of the registry modification. Look for anomalous behavior from the associated image and any unusual command-line arguments.
  - **Review Historical Changes**: Query the Windows registry history if available, examining previous values for the modified key to determine the legitimacy of changes made.
  - **Threat Intelligence Validation**: Cross-reference the detected modification against threat intelligence reports to identify whether similar activities are noted for any known threats or vulnerabilities.
  ### Prioritization
  This alert is classified as high severity due to the critical nature of registry modifications that can enable persistence for malware. Any alteration in essential run keys should be investigated promptly as it may indicate ongoing compromise or preparation for further malicious actions.
  ### Blind spots and Assumptions
  This rule may fail to fire if registry auditing is not properly configured or enabled on the monitored systems. Additionally, if an adversary uses legitimate administrative accounts or alters the behavior of legitimate administrative tools to disguise their activity, the changes may go undetected. It's also important to note that other legitimate applications may modify the same registry keys, which could lead to false positives if context is not adequately assessed.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
