# sigma/rules/windows/registry/registry_event/registry_event_persistence_recycle_bin.yml
title: Registry Persistence Mechanisms in Recycle Bin
id: 277efb8f-60be-4f10-b4d3-037802f37167
status: test
description: Detects persistence registry keys for Recycle Bin
references:
- https://github.com/vxunderground/VXUG-Papers/blob/751edb8d50f95bd7baa730adf2c6c3bb1b034276/The%20Persistence%20Series/Persistence%20via%20Recycle%20Bin/Persistence_via_Recycle_Bin.pdf
- https://persistence-info.github.io/Data/recyclebin.html
- https://www.hexacorn.com/blog/2018/05/28/beyond-good-ol-run-key-part-78-2/
author: frack113
date: 2021-11-18
modified: 2022-12-06
tags:
- attack.persistence
- attack.t1547
logsource:
  category: registry_event
  product: windows
detection:
  selection_create:
    EventType: RenameKey
    NewName|contains: '\CLSID\{645FF040-5081-101B-9F08-00AA002F954E}\shell\open'
  selection_set:
    EventType: SetValue
    TargetObject|contains: '\CLSID\{645FF040-5081-101B-9F08-00AA002F954E}\shell\open\command\(Default)'
  condition: 1 of selection_*
falsepositives:
- Unknown
level: high
notes: |
  ### Technical Context
  This Sigma rule is designed to detect persistence mechanisms within the Windows Registry that utilize the Recycle Bin. Specifically, it tracks alterations to specific registry keys associated with the Recycle Bin's CLSID (Class Identifier) that malware may leverage to establish persistence on a system. The rule identifies two event types: `RenameKey` for when registry keys are renamed, and `SetValue` for changes to values within a key. By monitoring these actions, it aims to uncover potential malicious activities that utilize the Recycle Bin for persistence, a technique documented under the MITRE ATT&CK framework as **T1547 
  - Boot or Logon Autostart Execution**. The primary log source for this detection is Windows Registry event logs, focusing on modifications that can indicate nefarious intent.
  ### Investigation Steps
  - **Review EDR Alerts:** Check the EDR tool for the alert triggering this rule to identify any associated processes or behaviors following the registry changes.
  - **Correlate with Process Creation Logs:** Utilize process creation logs to find any processes that were launched immediately before or after the registry modifications, which may indicate exploitation attempts.
  - **Analyze Command-Line Parameters:** Investigate command-line parameters of any related processes, specifically looking for unusual flags or indicators of tampering, especially those that leverage `%SystemRoot%\System32\recycler`.
  - **Check for Additional Persistence Mechanisms:** Look into other persistence mechanisms that may be linked to the affected system, such as startup folders or scheduled tasks, to assess the full scope of the potential compromise.
  ### Prioritization
  Given the alert is classified as high severity, it indicates a significant potential risk to the enterprise, suggesting that the persistence method may allow an adversary long-term access to the system, warranting immediate investigation and response.
  ### Blind Spots and Assumptions
  This rule may not capture all instances of malicious persistence techniques, especially if an adversary uses alternative methods to conceal their actions or employs legitimate tools excessively. Assumptions include that the registry events will be logged consistently and that the system in question is not in a state of normal administrative activity that might generate similar events as false positives (e.g., legitimate software installations or configurations). Additionally, if an attacker deletes or obfuscates the registry key before the monitoring tools can capture the activity, this rule may fail to fire.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
