# sigma/rules/windows/registry/registry_event/registry_event_susp_download_run_key.yml
title: Suspicious Run Key from Download
id: 9c5037d1-c568-49b3-88c7-9846a5bdc2be
status: test
description: Detects the suspicious RUN keys created by software located in Download
  or temporary Outlook/Internet Explorer directories
references:
- https://app.any.run/tasks/c5bef5b7-f484-4c43-9cf3-d5c5c7839def/
author: Florian Roth (Nextron Systems)
date: 2019-10-01
modified: 2021-11-27
tags:
- attack.persistence
- attack.t1547.001
logsource:
  category: registry_event
  product: windows
detection:
  selection:
    Image|contains:
    - '\Downloads\'
    - '\Temporary Internet Files\Content.Outlook\'
    - '\Local Settings\Temporary Internet Files\'
    TargetObject|contains: '\SOFTWARE\Microsoft\Windows\CurrentVersion\Run\'
  condition: selection
falsepositives:
- Software installers downloaded and used by users
level: high
notes: |
  n
  ### Technical Context
  This detection rule identifies potentially malicious activity by monitoring the creation of RUN keys in the Windows registry, specifically targeting keys created by software located in the Downloads folder or temporary directories associated with Outlook and Internet Explorer. The rule specifically looks for entries in the registry path `SOFTWARE\Microsoft\Windows\CurrentVersion\Run`, which is a common persistence mechanism leveraged by attackers to maintain control over infected systems. The analyzed data sources include Windows registry event logs and focus on the paths that are indicative of downloaded or temporary files, which may be used to execute malicious payloads automatically upon system startup. This behavior is linked to the MITRE ATT&CK framework under the tactic of Persistence and specifically the technique T1547.001, where attackers create or manipulate registry keys to ensure their payloads run with each system boot.
  ### Investigation Steps
  - **Review Registry Entries:** Examine the newly created RUN keys in the Windows registry to determine the provenance of the associated software. Use EDR tools to find the state before and after the creation to discern any unauthorized modifications.
    
  - **Analyze Process Creation Logs:** Utilize Sysmon or Windows Event Logs to trace the process creation events linked to the suspicious RUN keys. Pay attention to the command-line arguments and parent-child relationships of the processes spawned.
  - **Check for Suspicious Network Activity:** Analyze network connection logs using a Network Detection and Response (NDR) solution to look for unusual outbound connections that may correlate with the identified RUN keys. Focus on connections to known malicious IP addresses or domains.
  - **Inspect File Activity:** Investigate the files executed as a result of the RUN keys, including those in the Downloads and Temporary Internet Files directories. Use AV or EDR solutions to scan these files for known malware signatures or suspicious behavior.
  ### Prioritization
  Given that the alert is classified as high severity, it indicates that the creation of suspicious RUN keys may signify an ongoing compromise or malware installation attempt, which poses a significant risk to the organization's security posture.
  ### Blind spots and Assumptions
  The rule may not fire if the adversary uses a legitimate application to create the RUN keys or if the software was installed previously without detection. Furthermore, false positives may arise from legitimate software installers that create RUN keys during normal operation. This highlights an assumption that any entry in these paths is immediately suspect, which may not always be the case, necessitating context-specific investigation by analysts.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
