# sigma/rules/windows/builtin/taskscheduler/win_taskscheduler_execution_from_susp_locations.yml
title: Scheduled Task Executed From A Suspicious Location
id: 424273ea-7cf8-43a6-b712-375f925e481f
status: test
description: Detects the execution of Scheduled Tasks where the Program being run
  is located in a suspicious location or it's an unusale program to be run from a
  Scheduled Task
references:
- Internal Research
author: Nasreddine Bencherchali (Nextron Systems)
date: 2022-12-05
modified: 2023-02-07
tags:
- attack.persistence
- attack.t1053.005
logsource:
  product: windows
  service: taskscheduler
  definition: 'Requirements: The "Microsoft-Windows-TaskScheduler/Operational" is
    disabled by default and needs to be enabled in order for this detection to trigger'
detection:
  selection:
    EventID: 129     # Created Task Process
    Path|contains:
    - 'C:\Windows\Temp\'
    - '\AppData\Local\Temp\'
    - '\Desktop\'
    - '\Downloads\'
    - '\Users\Public\'
    - 'C:\Temp\'
    # If you experience FP. Uncomment the filter below and add the specific TaskName with the Program to it
    # filter:
    #     TaskName: '\Exact\Task\Name'
    #     Path: 'Exact\Path'
  condition: selection
falsepositives:
- Unknown
level: medium
notes: |
  ### Technical Context
  This Sigma rule is designed to detect the execution of Windows Scheduled Tasks that originate from suspicious directories such as temporary folders or user-specific locations that are not typical for task execution. The primary focus is on detecting Event ID 129 from the Task Scheduler operational log, which indicates that a task has successfully started. The program path is checked for known suspicious locations, including common directories like `C:\Windows\Temp\`, `\AppData\Local\Temp\`, `\Downloads\`, and other non-standard locations. The intention behind this detection is to identify potentially malicious activities, as attackers often leverage Scheduled Tasks to maintain persistence using executables stored in these less scrutinized directories. This rule corresponds to the MITRE ATT&CK tactic of **Persistence** and the technique **T1053.005** (Scheduled Task/Job: Scheduled Task Created by User).
  ### Investigation Steps
  - **Check Task Scheduler Logs**: Review the Task Scheduler logs associated with Event ID 129 to identify the specific task that was executed, including the user context and time of execution.
    
  - **Correlate with Process Creation Logs**: Utilize EDR tools to correlate the execution of the Scheduled Task with recent process creation events, ensuring to review the command-line arguments provided at runtime.
  - **Review File Paths and Executables**: Inspect the executable being run by the task for any known malicious indicators or unusual behavior and verify its file hash against threat intelligence sources.
  - **Examine Network Activity**: Utilize network logs to assess any outbound connections made by the application initiated through the Scheduled Task, which can indicate communication with command and control servers.
  ### Prioritization
  This alert should be prioritized as medium severity due to the potential for malicious persistence mechanisms that can be exploited by threat actors. Investigating these alerts is crucial in establishing whether the Scheduled Task represents an active threat or a benign administrative task.
  ### Blind Spots and Assumptions
  The detection may fail to trigger if the Task Scheduler operational logs are not enabled, as they are disabled by default in Windows environments. Additionally, legitimate applications may be mistakenly flagged if they operate from the specified paths. Attackers may also employ techniques to obfuscate their activities, such as using more legitimate paths or legitimate-looking executables, hence bypassing this detection entirely. Analysts should remain aware that not all potentially malicious tasks will originate from these predefined paths.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality and usefulness, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
