# sigma/rules/windows/file/file_event/file_event_win_office_outlook_susp_macro_creation.yml
title: Suspicious Outlook Macro Created
id: 117d3d3a-755c-4a61-b23e-9171146d094c
related:
- id: 8c31f563-f9a7-450c-bfa8-35f8f32f1f61
  type: derived
status: test
description: Detects the creation of a macro file for Outlook.
references:
- https://www.mdsec.co.uk/2020/11/a-fresh-outlook-on-mail-based-persistence/
- https://speakerdeck.com/heirhabarov/hunting-for-persistence-via-microsoft-exchange-server-or-outlook?slide=53
- https://www.linkedin.com/pulse/outlook-backdoor-using-vba-samir-b-/
author: Nasreddine Bencherchali (Nextron Systems)
date: 2023-02-08
tags:
- attack.persistence
- attack.command-and-control
- attack.t1137
- attack.t1008
- attack.t1546
logsource:
  category: file_event
  product: windows
detection:
  selection:
    TargetFilename|endswith: '\Microsoft\Outlook\VbaProject.OTM'
  filter:
    Image|endswith: '\outlook.exe'
  condition: selection and not filter
falsepositives:
- Unlikely
level: high
notes: |
  ### Technical Context
  This Sigma rule is designed to detect the creation of a macro file specifically for Microsoft Outlook, which can be a tactic used by threat actors to establish persistence and conduct command-and-control activities. It focuses on file events generated when the `VbaProject.OTM` file—commonly associated with Outlook macros—is created or modified. The rule examines the `TargetFilename` field to identify files that end with `\Microsoft\Outlook\VbaProject.OTM`, indicating a potentially malicious macro creation. Additionally, it filters out events where the process responsible for the file creation is not `outlook.exe`, ensuring that the detection focuses on legitimate Outlook macro activity. This detection aligns with the MITRE ATT&CK techniques T1137 (Human Interface Device) and T1546 (Event Triggered Execution), which encompass processes that may be exploited for malicious purposes.
  ### Investigation Steps
  - **Correlate Process Activity:** Review the process creation logs to determine if `outlook.exe` was executing at the time of the suspect macro creation. Check if there are any unusual command-line parameters associated with the Outlook process.
  - **Examine Outlook Macro Environment:** Investigate the user's Outlook environment, including reviewing existing macros, add-ins, and any scheduled tasks that may interact with the macro framework. Look for unusual or unauthorized macros that may indicate persistence mechanisms.
  - **Analyze Email Behavior:** Review the user's recent email activity to see if they have interacted with any suspicious or unsolicited emails that may have contained malicious payloads capable of triggering the macro's creation.
  - **Check Security Events:** Utilize your Security Information and Event Management (SIEM) system to analyze related security events and alerts that occurred in conjunction with the timestamp of the detected macro creation. This could help identify any prior or subsequent suspicious activities.
  ### Prioritization
  This alert is considered high severity due to the potential for an attacker to gain continued access to the environment through malicious macros, often leading to further exploitation or data exfiltration.
  ### Blind Spots and Assumptions
  It is important to note that well-crafted attacks may avoid detection by either modifying the file paths or using legitimate user credentials. Additionally, users may create benign macros as part of their legitimate workflow, leading to potential false positives. Furthermore, the rule relies on the correct logging configuration; if events are not logged or there are gaps in logging strategies, the rule may not fire. Finally, adversaries might leverage legitimate processes or applications, making the detection of malicious activity more challenging.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
