# sigma/rules/windows/file/file_event/file_event_win_susp_startup_folder_persistence.yml
title: Suspicious Startup Folder Persistence
id: 28208707-fe31-437f-9a7f-4b1108b94d2e
related:
- id: 2aa0a6b4-a865-495b-ab51-c28249537b75
  type: similar
status: test
description: Detects when a file with a suspicious extension is created in the startup
  folder
references:
- https://github.com/last-byte/PersistenceSniper
author: Nasreddine Bencherchali (Nextron Systems)
date: 2022-08-10
modified: 2023-01-06
tags:
- attack.persistence
- attack.t1547.001
logsource:
  product: windows
  category: file_event
detection:
  selection:
    TargetFilename|contains: '\Windows\Start Menu\Programs\Startup\'
    TargetFilename|endswith:
            # Add or remove suspicious extensions according to your env needs
    - '.vbs'
    - '.vbe'
    - '.bat'
    - '.ps1'
    - '.hta'
    - '.dll'
    - '.jar'
    - '.msi'
    - '.scr'
    - '.cmd'
  condition: selection
falsepositives:
- Rare legitimate usage of some of the extensions mentioned in the rule
level: high
notes: |
  ### Technical Context
  The "Suspicious Startup Folder Persistence" Sigma rule is designed to detect the creation of files with specified suspicious extensions within the Windows Startup folder. This folder serves as a common place for ensuring programs automatically launch when a user logs into their account. The rule specifically monitors file creations, focusing on extensions commonly associated with scripting, such as `.vbs`, `.ps1`, and executable file types including `.msi` and `.dll`. By leveraging file event logs, the rule assesses whether any new entries in the Startup directory match the criteria, indicating a potential persistence mechanism employed by adversaries as part of the MITRE ATT&CK technique for "Boot or Logon Autostart Execution" (T1547.001). 
  The primary data source involved in this detection is the Windows File Event logs, which track file-related activities, enabling the SOC to identify unauthorized or potentially malicious modifications to the startup behavior on systems.
  ### Investigation Steps
  - **Assess the Alert Context:** Review the alert details in the SIEM for the file path, filename, and timestamp. Cross-reference these with known legitimate applications to determine the validity of the file creation.
  - **Correlate EDR Data:** Utilize EDR tools to examine the process creation logs around the time the suspicious file was created. Look for abnormal command-line parameters or parent-child process relationships that may indicate malicious activity.
  - **Check for User Activity:** Investigate the user account that generated the file using Windows Security logs. Determine if the user's account shows signs of compromise or if they were acting anomalously.
  - **Analyze the Startup Folder:** Directly review the contents of the Startup folder on the affected system. Assess whether any additional suspicious files exist, or if modifications were made to legitimate entries that could indicate tampering.
  ### Prioritization
  This alert is categorized with a high severity due to the potential for persistence mechanisms that allow adversaries to maintain access to compromised systems, making it critical for timely investigation to prevent further exploitation.
  ### Blind Spots and Assumptions
  This rule assumes that the monitoring of file event logs is properly configured and that the suspicious file extensions list includes all potential malicious extensions applicable to the environment. There may be blind spots related to variations in legitimate software behavior that use the same extensions, which could lead to false positives. Moreover, sophisticated adversaries may employ techniques to disguise their payloads or utilize legitimate files, bypassing this detection.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
