# sigma/rules/windows/file/file_event/file_event_win_msdt_susp_directories.yml
title: File Creation In Suspicious Directory By Msdt.EXE
id: 318557a5-150c-4c8d-b70e-a9910e199857
status: test
description: Detects msdt.exe creating files in suspicious directories which could
  be a sign of exploitation of either Follina or Dogwalk vulnerabilities
references:
- https://irsl.medium.com/the-trouble-with-microsofts-troubleshooters-6e32fc80b8bd
- https://msrc-blog.microsoft.com/2022/05/30/guidance-for-cve-2022-30190-microsoft-support-diagnostic-tool-vulnerability/
author: Vadim Varganov, Florian Roth (Nextron Systems)
date: 2022-08-24
modified: 2023-02-23
tags:
- attack.persistence
- attack.t1547.001
- cve.2022-30190
logsource:
  category: file_event
  product: windows
detection:
  selection:
    Image|endswith: '\msdt.exe'
    TargetFilename|contains:
    - '\Desktop\'
    - '\Start Menu\Programs\Startup\'
    - 'C:\PerfLogs\'
    - 'C:\ProgramData\'
    - 'C:\Users\Public\'
  condition: selection
falsepositives:
- Unknown
level: high
notes: |
  ### Technical Context
  The Sigma rule titled "File Creation In Suspicious Directory By Msdt.EXE" is designed to detect the creation of files by the `msdt.exe` process in specific directories that are commonly targeted for malicious activity. This behavior is often indicative of exploitation attempts relating to the Follina vulnerability (CVE-2022-30190) and Dogwalk vulnerabilities, where attackers utilize legitimate system utilities to execute harmful actions without raising suspicion. The rule specifically monitors file creation events in well-known locations such as the Desktop, Startup folder, ProgramData, and other paths that can be exploited for persistence or execution. Key technical data sources involved in this detection include Windows Security logs that capture file events when `msdt.exe` attempts to create files in these directories. This aligns with the MITRE ATT&CK framework, specifically the tactics of Persistence and the technique T1547.001 (Boot or Logon Autostart Execution).
  ### Investigation Steps
  - **Verify File Creation:** Review recent file creation events captured in Windows Security logs to confirm whether files were created in the targeted suspicious directories by `msdt.exe`.
  - **Check Process Context:** Utilize EDR tools to analyze the context of the `msdt.exe` process, including any associated parent processes and command-line arguments that may provide insights into its execution parameters.
  - **Assess Compromised Systems:** Investigate the endpoints where this alert was triggered to check for other indicators of compromise (IOCs) or additional malicious activity potentially associated with the exploitation of the vulnerabilities.
  - **Review Additional Logs:** Correlate findings with other telemetry sources, such as Sysmon logs, to identify any abnormal network connections or suspicious command executions linked to the timing of the file creation events.
  ### Prioritization
  This alert is marked as high severity due to its association with known exploitation techniques that can lead to persistent malware installation, making immediate investigation crucial to mitigate potential threats within the enterprise environment.
  ### Blind Spots and Assumptions
  This rule may not fire under certain conditions, such as if the `msdt.exe` process runs in a non-standard way or if an attacker employs other methods to obscure their activities. Additionally, there may be benign file creations in the monitored directories that could lead to false positives. Engineers should be aware that adversaries may use legitimate processes for malicious purposes, complicating detection efforts.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
