# sigma/rules/windows/file/file_event/file_event_win_ripzip_attack.yml
title: Potential RipZip Attack on Startup Folder
id: a6976974-ea6f-4e97-818e-ea08625c52cb
status: test
description: |
  Detects a phishing attack which expands a ZIP file containing a malicious shortcut.
  If the victim expands the ZIP file via the explorer process, then the explorer process expands the malicious ZIP file and drops a malicious shortcut redirected to a backdoor into the Startup folder.
  Additionally, the file name of the malicious shortcut in Startup folder contains {0AFACED1-E828-11D1-9187-B532F1E9575D} meaning the folder shortcut operation.
references:
- https://twitter.com/jonasLyk/status/1549338335243534336?t=CrmPocBGLbDyE4p6zTX1cg&s=19
author: Greg (rule)
date: 2022-07-21
modified: 2023-01-05
tags:
- attack.persistence
- attack.t1547
logsource:
  category: file_event
  product: windows
detection:
  selection:   # %APPDATA%\Microsoft\Windows\Start Menu\Programs\Startup\target.lnk.{0AFACED1-E828-11D1-9187-B532F1E9575D}\target.lnk
    TargetFilename|contains|all:
    - '\Microsoft\Windows\Start Menu\Programs\Startup'
    - '.lnk.{0AFACED1-E828-11D1-9187-B532F1E9575D}'
    Image|endswith: '\explorer.exe'
  condition: selection
falsepositives:
- Unknown
level: high
notes: |
  n
  ### Technical Context
  The Sigma rule titled "Potential RipZip Attack on Startup Folder" identifies a specific persistence technique that exploits user actions to compromise a system. This attack vector involves a phishing attempt where a ZIP file—when extracted by the victim using the Windows Explorer process—drops a malicious shortcut file in the Startup folder. The dropped shortcut is unique as its filename contains a specific GUID: `{0AFACED1-E828-11D1-9187-B532F1E9575D}`. Such a signature indicates that this shortcut is intended to execute during the startup sequence, thereby maintaining the malicious activity each time the system boots. The rule monitors Windows file events, particularly focusing on creation activities in the Startup folder, with the intention of detecting this malicious LNK file tied to the malicious exploit. This technique corresponds to the MITRE ATT&CK persistence tactic (T1547), which encompasses methods adversaries use to maintain a foothold on a compromised system.
  ### Investigation Steps
  - **Verify the Malicious Shortcut:** Check the contents and properties of any `.lnk` files located in the `%APPDATA%\Microsoft\Windows\Start Menu\Programs\Startup\` directory, specifically those containing the GUID `{0AFACED1-E828-11D1-9187-B532F1E9575D}`.
  - **Analyze Process Activity:** Utilize EDR tools to review process creation logs around the time the malicious shortcut was created, focusing on instances of `explorer.exe` to identify the extraction of the ZIP file.
  - **Inspect Network Connections:** Use network detection tools to track outbound connections initiated by the malicious shortcut, which may be indicative of a backdoor attempting to communicate with a remote server.
  - **Examine User Activity:** Look into Windows Event Logs for any suspicious user activities or logs that might show the user's interaction with the extraction process, particularly around the time of the initial alert trigger.
  ### Prioritization
  This alert is categorized as high severity due to the potential for a backdoor to be established on the endpoint, allowing an adversary to maintain ongoing access and control over the compromised system.
  ### Blind spots and Assumptions
  This rule may not trigger in environments where users typically access ZIP files using tools other than Windows Explorer, or if the malicious files are deleted before detection. Additionally, the existence of security solutions (like AV software) that might quarantine or delete such files could prevent these activities from being logged. Adversaries could also use obfuscation techniques or change the GUID used in the LNK file to avoid detection by this specific rule.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality and validity, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
