# sigma/rules/windows/file/file_event/file_event_win_susp_task_write.yml
title: Suspicious Scheduled Task Write to System32 Tasks
id: 80e1f67a-4596-4351-98f5-a9c3efabac95
status: test
description: Detects the creation of tasks from processes executed from suspicious
  locations
references:
- Internal Research
author: Florian Roth (Nextron Systems)
date: 2021-11-16
modified: 2022-01-12
tags:
- attack.persistence
- attack.execution
- attack.t1053
logsource:
  product: windows
  category: file_event
detection:
  selection:
    TargetFilename|contains: '\Windows\System32\Tasks'
    Image|contains:
    - '\AppData\'
    - 'C:\PerfLogs'
    - '\Windows\System32\config\systemprofile'
  condition: selection
falsepositives:
- Unknown
level: high
notes: |
  ### Technical Context
  This detection rule aims to identify the creation of scheduled tasks within the Windows System32 Tasks directory, particularly from processes originating from suspicious locations such as user AppData folders, performance logs, or the system profile directory. Task creation in System32 can often signal malicious intent, as attackers may aim to persist their malware or perform scheduled actions to automate further exploitation. The rule monitors file events in Windows, focusing on specific characteristics that indicate potentially illicit behavior. The relevant MITRE ATT&CK technique for this detection is **T1053: Scheduled Task/Job**, which falls under the persistence and execution tactics. By correlating process origins with task creation activities, the rule seeks to surface instances of potentially harmful behaviors that warrant further investigation.
  ### Investigation Steps
  - **Review Task Creation Events:** Examine the event logs for any task creation events pointing towards suspicious activities. Focus on the user accounts associated with these tasks to assess if they are legitimate or potentially compromised.
  - **Analyze Process Origins:** Utilize EDR tools to investigate the processes that triggered the task creation, particularly their origins and command-line parameters. Confirm if these processes are indeed executing from uncommon directories like AppData or system profile.
  - **Cross-reference with Threat Intelligence:** Check the identified processes against threat intelligence repositories to ascertain if they are known or suspected malicious applications or behaviors. This can provide additional context regarding any potential threat actors.
  - **Monitor for Subsequent Activity:** After identifying suspicious task creation, maintain surveillance on the affected systems for further indicators of compromise, such as unusual network connections, additional process creations, or unauthorized file access.
  ### Prioritization
  The alert is classified as high severity because the creation of scheduled tasks from suspicious locations is commonly indicative of persistent threats or malicious behavior, which can lead to significant security risks and operational disruptions if not addressed promptly.
  ### Blind Spots and Assumptions
  This rule may not fire in scenarios where the task creation is executed from clean processes or directories that have been whitelisted. Furthermore, if an adversary is using more sophisticated techniques to hide their activities, like running processes that evade detection or employing legitimate applications under the guise of normalcy, this detection could be bypassed. It assumes that all scheduled tasks will generate logs that can be monitored, which may not be the case if logging is improperly configured or if the settings have been modified by an adversary.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
