# sigma/rules/windows/registry/registry_set/registry_set_vbs_payload_stored.yml
title: VBScript Payload Stored in Registry
id: 46490193-1b22-4c29-bdd6-5bf63907216f
status: test
description: Detects VBScript content stored into registry keys as seen being used
  by UNC2452 group
references:
- https://www.microsoft.com/security/blog/2021/03/04/goldmax-goldfinder-sibot-analyzing-nobelium-malware/
author: Florian Roth (Nextron Systems)
date: 2021-03-05
modified: 2023-08-17
tags:
- attack.persistence
- attack.t1547.001
logsource:
  category: registry_set
  product: windows
detection:
  selection:
    TargetObject|contains: 'Software\Microsoft\Windows\CurrentVersion'
    Details|contains:
    - 'vbscript:'
    - 'jscript:'
    - 'mshtml,'
    - 'RunHTMLApplication'
    - 'Execute('
    - 'CreateObject'
    - 'window.close'
  filter:
    TargetObject|contains: 'Software\Microsoft\Windows\CurrentVersion\Run'
  filter_dotnet:
    Image|endswith: '\msiexec.exe'
    TargetObject|contains: '\SOFTWARE\Microsoft\Windows\CurrentVersion\Installer\UserData\'
    Details|contains:
    - '\Microsoft.NET\Primary Interop Assemblies\Microsoft.mshtml.dll'
    - '<\Microsoft.mshtml,fileVersion='
    - '_mshtml_dll_'
    - '<\Microsoft.mshtml,culture='
  condition: selection and not 1 of filter*
falsepositives:
- Unknown
level: high
notes: |
  ### Technical Context
  The Sigma rule "VBScript Payload Stored in Registry" is designed to detect suspicious entries in the Windows Registry that may indicate the presence of malicious VBScript content. Attackers, such as the UNC2452 group, often utilize techniques that persistently store payloads in registry keys to maintain control over compromised systems. This rule focuses on detecting alterations in the relevant registry paths—specifically, those under `Software\Microsoft\Windows\CurrentVersion`, which is commonly used for executing scripts at startup. 
  The criteria for triggering this alert include the presence of particular strings associated with VBScript execution (e.g., 'vbscript:', 'RunHTMLApplication', 'CreateObject') within the details of registry modifications. The rule is aligned with the MITRE ATT&CK framework under the tactic of Persistence and the technique T1547.001, indicating a method of gaining persistent access through registry manipulation.
  ### Investigation Steps
  - **Review the Alert Detail**: Start by examining the specific registry keys that triggered the alert within the EDR tool. Look for indicators of compromise such as the presence of suspicious script commands or unusual user context.
  - **Correlate with Process Activity**: Utilize Sysmon or Windows Security logs to list processes that were running around the time of the registry modification. Pay particular attention to processes that may have executed scripts or unusual command lines.
  - **Check Recent Changes**: Investigate any recent registry changes leading up to the alert. This can help determine if legitimate updates were occurring, or if the modification was malicious in nature.
  - **Cross-Reference with External Threat Intelligence**: Use threat intelligence platforms to analyze any identified registry keys and payloads against known bad actors or malware families, verifying if there are documented associations or behaviors.
  ### Prioritization
  This alert is classified as high severity due to the potential for persistent malware that can compromise system integrity and maintain unauthorized access. Since such behaviors are often indicative of advanced persistent threats, prompt investigation is crucial to prevent further exploitation.
  ### Blind Spots and Assumptions
  This detection rule may not fire in scenarios where the attacker bypasses registry modifications, such as utilizing alternative persistence mechanisms (e.g., scheduled tasks or WMI). Additionally, if the registry changes are made by legitimate applications that leverage scripts for functionality, false positives may occur. The accuracy of detection can also be impacted by logging configurations; if logging for certain registry actions is not enabled or incomplete, critical events may go unnoticed. 
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality and relevance, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
