# sigma/rules/windows/registry/registry_set/registry_set_netsh_help_dll_persistence_susp_location.yml
title: New Netsh Helper DLL Registered From A Suspicious Location
id: e7b18879-676e-4a0e-ae18-27039185a8e7
related:
- id: 56321594-9087-49d9-bf10-524fe8479452
  type: similar
- id: c90362e0-2df3-4e61-94fe-b37615814cb1
  type: similar
status: test
description: |
  Detects changes to the Netsh registry key to add a new DLL value that is located on a suspicious location. This change might be an indication of a potential persistence attempt by adding a malicious Netsh helper
references:
- https://www.ired.team/offensive-security/persistence/t1128-netsh-helper-dll
- https://pentestlab.blog/2019/10/29/persistence-netsh-helper-dll/
author: Nasreddine Bencherchali (Nextron Systems)
date: 2023-11-28
tags:
- attack.persistence
- attack.t1546.007
logsource:
  category: registry_set
  product: windows
detection:
  selection_target:
    TargetObject|contains: '\SOFTWARE\Microsoft\NetSh'
  selection_folders_1:
    Details|contains:
    - ':\Perflogs\'
    - ':\Users\Public\'
    - ':\Windows\Temp\'
    - '\AppData\Local\Temp\'
    - '\Temporary Internet'
  selection_folders_2:
  - Details|contains|all:
    - ':\Users\'
    - '\Favorites\'
  - Details|contains|all:
    - ':\Users\'
    - '\Favourites\'
  - Details|contains|all:
    - ':\Users\'
    - '\Contacts\'
  - Details|contains|all:
    - ':\Users\'
    - '\Pictures\'
  condition: selection_target and 1 of selection_folders_*
falsepositives:
- Unknown
level: high
notes: |
  n
  ### Technical Context
  This Sigma rule is designed to detect suspicious modifications to the Windows registry, specifically targeting the Netsh helper DLL registration. By keeping an eye on changes to the registry key located at `HKEY_LOCAL_MACHINE\SOFTWARE\Microsoft\NetSh`, the rule identifies potential persistence attempts by malicious actors who might register a new DLL from uncommon directories, such as `C:\Users\Public` or `C:\Windows\Temp`. These locations are often leveraged by attackers to establish a foothold within an environment without raising immediate red flags, as they are not typical locations for legitimate applications. The underlying data sources for this detection are Windows registry change events, specifically focusing on keys associated with networking and system configurations. The detection aligns with the MITRE ATT&CK technique for persistence (T1546.007), emphasizing the importance of scrutinizing abnormal registry activity for signs of compromise.
  ### Investigation Steps
  - **Check Registry Changes**: Use the Windows Registry editor or a SIEM tool to review the changes made to the `HKEY_LOCAL_MACHINE\SOFTWARE\Microsoft\NetSh` key, focusing on any new DLL entries added.
  - **Identify DLL Source**: Investigate the file path of the registered DLL to determine if it resides in a suspicious or uncommon directory, cross-referencing with known good paths in your environment.
  - **Review Process Creation Logs**: Utilize EDR tools to examine process creation logs corresponding to the time of the registry change, looking for any unusual process activity that corresponds with the DLL registration.
  - **Correlate Network Activity**: Analyze network connection logs, particularly focusing on any outbound connections initiated by processes associated with the suspicious DLL to detect potential exfiltration or callback attempts.
  ### Prioritization
  Given that the alert indicates a high severity level, this suggests a significant risk to the organization's security posture, as it may point to an attempt by an adversary to establish persistence and maintain access within the environment.
  ### Blind Spots and Assumptions
  This detection rule may not fire in situations where an adversary uses legitimate software tools or modifies the registry in a less overt manner, such as through scheduled tasks or other less common persistence techniques. Additionally, the rule assumes that proper logging is enabled for registry changes and that no log tampering has occurred. There may also be false positives if legitimate applications are installed from the identified directories, highlighting the need for extensive context during investigations.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
