# sigma/rules/windows/registry/registry_set/registry_set_taskcache_entry.yml
title: Scheduled TaskCache Change by Uncommon Program
id: 4720b7df-40c3-48fd-bbdf-fd4b3c464f0d
status: test
description: Monitor the creation of a new key under 'TaskCache' when a new scheduled
  task is registered by a process that is not svchost.exe, which is suspicious
references:
- https://thedfirreport.com/2021/03/29/sodinokibi-aka-revil-ransomware/
- https://labs.f-secure.com/blog/scheduled-task-tampering/
author: Syed Hasan (@syedhasan009)
date: 2021-06-18
modified: 2023-08-17
tags:
- attack.persistence
- attack.t1053
- attack.t1053.005
logsource:
  category: registry_set
  product: windows
detection:
  selection:
    TargetObject|contains: 'SOFTWARE\Microsoft\Windows NT\CurrentVersion\Schedule\TaskCache\'
  filter:
    TargetObject|contains:
    - 'Microsoft\Windows\UpdateOrchestrator'
    - 'Microsoft\Windows\SoftwareProtectionPlatform\SvcRestartTask\Index'
    - 'Microsoft\Windows\Flighting\OneSettings\RefreshCache\Index'
  filter_tiworker:
    Image|startswith: 'C:\Windows\'
    Image|endswith: '\TiWorker.exe'
  filter_svchost:
    Image: 'C:\WINDOWS\system32\svchost.exe'
  filter_ngen:
    Image|startswith: 'C:\Windows\Microsoft.NET\Framework'     # \Framework\ and \Framework64\
    Image|endswith: '\ngen.exe'
    TargetObject|contains:
    - '\Microsoft\Windows NT\CurrentVersion\Schedule\TaskCache\Tasks\{B66B135D-DA06-4FC4-95F8-7458E1D10129}'
    - '\Microsoft\Windows NT\CurrentVersion\Schedule\TaskCache\Tree\Microsoft\Windows\.NET
      Framework\.NET Framework NGEN'
  filter_office_click_to_run:
    Image:
    - 'C:\Program Files\Microsoft Office\root\Integration\Integrator.exe'
    - 'C:\Program Files (x86)\Microsoft Office\root\Integration\Integrator.exe'
  filter_msiexec:
    Image: 'C:\Windows\System32\msiexec.exe'
  filter_dropbox_updater:
    Image:
    - 'C:\Program Files (x86)\Dropbox\Update\DropboxUpdate.exe'
    - 'C:\Program Files\Dropbox\Update\DropboxUpdate.exe'
  filter_explorer:
    Image: 'C:\Windows\explorer.exe'
    TargetObject|contains: '\Microsoft\Windows NT\CurrentVersion\Schedule\TaskCache\Tree\Microsoft\Windows\PLA\Server
      Manager Performance Monitor\'
  filter_system:
    Image: 'System'
  condition: selection and not 1 of filter*
falsepositives:
- Unknown
level: high
notes: |
  ### Technical Context
  The "Scheduled TaskCache Change by Uncommon Program" rule monitors the Windows registry for the creation of new keys under `TaskCache`, specifically when a new scheduled task is registered by processes other than the commonly trusted `svchost.exe`. Scheduled tasks can be a legitimate mechanism for automating functions in Windows; however, adversaries often exploit this feature to maintain persistence or execute malicious payloads. The rule utilizes Windows registry logs as its primary data source, particularly focusing on registry operations that manipulate scheduled tasks. By checking for the `TaskCache` registry path, it captures potentially suspicious activity that might be indicative of unauthorized task creation, in alignment with the persistence tactics outlined in the MITRE ATT&CK framework, notably T1053 and T1053.005.
  ### Investigation Steps
  - **Verify Process Source:** Utilize EDR tools to investigate the process that created the registry key. Check its reputation and history for any known malicious activity or anomalies.
  - **Analyze Registry Changes:** Review the specific task that was registered by examining the newly created keys in `TaskCache`. Determine the legitimacy of the task by looking into its intended function and target systems.
  - **Inspect Other Logs:** Correlate with Sysmon logs to see if the unusual process executed any other suspicious actions, such as process creation or network connections around the same time.
  - **Check for Additional Indicators:** Query network logs (NDR/NGFW) for any outbound connections made by the identified process, as well as related user activity in proxy and VPN logs to further assess potential compromise.
  ### Prioritization
  The high severity level assigned to this alert is due to the potential for unauthorized modifications to scheduled tasks which can serve as a persistent backdoor for an attacker, making this a critical indicator of compromise that warrants immediate investigation.
  ### Blind Spots and Assumptions
  This rule may not fire in environments where uncommon programs legitimately create scheduled tasks, leading to potential blind spots in detection. Additionally, if an attacker uses a whitelisted application that meets the filter criteria, it could evade detection. Adversaries may also exploit legitimate system processes in forms that do not create clear registry indicators, thereby bypassing this particular rule. 
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
