# sigma/rules/windows/registry/registry_set/registry_set_powershell_in_run_keys.yml
title: Suspicious Powershell In Registry Run Keys
id: 8d85cf08-bf97-4260-ba49-986a2a65129c
status: test
description: Detects potential PowerShell commands or code within registry run keys
references:
- https://github.com/frack113/atomic-red-team/blob/a9051c38de8a5320b31c7039efcbd3b56cf2d65a/atomics/T1547.001/T1547.001.md#atomic-test-9---systembc-malware-as-a-service-registry
- https://www.trendmicro.com/en_us/research/22/j/lv-ransomware-exploits-proxyshell-in-attack.html
author: frack113, Florian Roth (Nextron Systems)
date: 2022-03-17
modified: 2023-08-17
tags:
- attack.persistence
- attack.t1547.001
logsource:
  category: registry_set
  product: windows
detection:
  selection:
    TargetObject|contains: '\Software\Microsoft\Windows\CurrentVersion\Run'     # Also covers "RunOnce" and "RunOnceEx"
    Details|contains:
    - 'powershell'
    - 'pwsh '
    - 'FromBase64String'
    - '.DownloadFile('
    - '.DownloadString('
    - ' -w hidden '
    - ' -w 1 '
    - '-windowstyle hidden'
    - '-window hidden'
    - ' -nop '
    - ' -encodedcommand '
    - '-ExecutionPolicy Bypass'
    - 'Invoke-Expression'
    - 'IEX ('
    - 'Invoke-Command'
    - 'ICM -'
    - 'Invoke-WebRequest'
    - 'IWR '
    - ' -noni '
    - ' -noninteractive '
  condition: selection
falsepositives:
- Legitimate admin or third party scripts. Baseline according to your environment
level: medium
notes: |
  ### Technical Context
  This Sigma rule detects potentially malicious PowerShell commands that may be used for persistence within the Windows Registry, specifically targeting the `Run` and `RunOnce` keys. By analyzing the registry changes made to these keys, the rule identifies variations of PowerShell commands or other suspicious command patterns that could indicate an attempt to execute malicious scripts or payloads upon system startup. Key indicators include common PowerShell flags and methods such as `-ExecutionPolicy Bypass` and `Invoke-WebRequest`, often used by attackers to execute scripts discreetly. This rule primarily leverages Windows Registry logs as the data source and corresponds to the MITRE ATT&CK tactic of persistence (T1547.001).
  ### Investigation Steps
  - **Check Registry Keys:** Investigate the contents of the `HKEY_CURRENT_USER\Software\Microsoft\Windows\CurrentVersion\Run` and `RunOnce` keys for unusual entries that involve PowerShell commands or encoded scripts.
  - **Analyze Process Creation:** Utilize EDR to track process creation logs following registry modifications and see if any suspicious processes were initiated in relation to the recorded registry changes.
  - **Inspect System Event Logs:** Review Windows Security and System event logs for any unexpected account logins or application installations that correlate with the timing of the registry modifications.
  - **Monitor Network Connections:** Use NDR tools to examine network traffic for any outbound connections initiated by the identified PowerShell commands, specifically looking for unusual domains or file downloads.
  ### Prioritization
  The alert is assigned a medium severity level due to the potential for establishing persistence on compromised systems, which may lead to further malicious actions, including data exfiltration or lateral movement within the network.
  ### Blind spots and Assumptions
  This rule may not fire if legitimate administrative scripts are used that do not match the specified PowerShell patterns, highlighting the need for environment-specific baselining. Additionally, malicious actors could obfuscate their commands or utilize less common persistence mechanisms, potentially evading detection. Engineers should recognize that sophisticated adversaries may adapt their techniques over time, making continuous adjustment and validation of the detection rule necessary.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
