# sigma/rules/windows/registry/registry_set/registry_set_susp_run_key_img_folder.yml
title: New RUN Key Pointing to Suspicious Folder
id: 02ee49e2-e294-4d0f-9278-f5b3212fc588
status: experimental
description: Detects suspicious new RUN key element pointing to an executable in a
  suspicious folder
references:
- https://www.fireeye.com/blog/threat-research/2018/08/fin7-pursuing-an-enigmatic-and-evasive-global-criminal-operation.html
author: Florian Roth (Nextron Systems), Markus Neis, Sander Wiebing
date: 2018-08-25
modified: 2024-07-16
tags:
- attack.persistence
- attack.t1547.001
logsource:
  category: registry_set
  product: windows
detection:
  selection_target:
    TargetObject|contains:
    - '\SOFTWARE\Microsoft\Windows\CurrentVersion\Run\'
    - '\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce\'
  selection_details:
  - Details|contains:
    - ':\$Recycle.bin\'
    - ':\Temp\'
    - ':\Users\Default\'
    - ':\Users\Desktop\'
    - ':\Users\Public\'
    - ':\Windows\Temp\'
    - '\AppData\Local\Temp\'
    - '%temp%\'
    - '%tmp%\'
  - Details|startswith:
    - '%Public%\'
    - 'wscript'
    - 'cscript'
  filter_main_windows_update:
    TargetObject|contains: '\Microsoft\Windows\CurrentVersion\RunOnce\'
    Image|startswith: 'C:\Windows\SoftwareDistribution\Download\'
    Details|contains|all:
    - 'rundll32.exe '
    - 'C:\WINDOWS\system32\advpack.dll,DelNodeRunDLL32'
    Details|contains:
    - '\AppData\Local\Temp\'
    - 'C:\Windows\Temp\'
  condition: all of selection_* and not 1 of filter_main_*
falsepositives:
- Software using weird folders for updates
level: high
notes: |
  n
  ### Technical Context
  This detection rule identifies the creation of new entries in the Windows registry's RUN and RUNOnce keys that point to executables located in suspicious and commonly abused directories such as Temp or Recycle Bin. By monitoring for these registry changes, the rule aims to flag potential persistence mechanisms that may indicate malicious activity according to the MITRE ATT&CK technique T1547.001 
  - Registry Run Keys / Startup Folder. The relevant telemetry includes Windows registry logs, specifically focusing on `TargetObject` that points to the mentioned paths, along with details that might indicate the path or properties of the executable.
  ### Investigation Steps
  - **Use EDR for Process Analysis:** Review the process tree in your EDR tool to understand what process created the suspicious RUN key entry and whether it invoked any unusual or unauthorized executables.
  - **Examine Registry Logs:** Search for recent registry modifications to `HKCU\Software\Microsoft\Windows\CurrentVersion\Run` and `HKCU\Software\Microsoft\Windows\CurrentVersion\RunOnce` using the relevant logs to identify other potential persistence mechanisms.
  - **Check File Locations:** Inspect the folder referenced in the RUN key entry for any suspicious files or executables, especially in the paths like Temp or Recycle Bin, and analyze their behavior with your antivirus or malware analysis tools.
  - **Network Traffic Review:** Investigate related network connections from the executing processes to check for any outbound communication to known malicious IP addresses or domains.
  ### Prioritization
  The alert is classified as high severity due to the potential association with persistent malware that can re-establish itself after a reboot, making it crucial for timely remediation efforts.
  ### Blind Spots and Assumptions
  This rule may not fire if malicious actors utilize legitimate processes or if they obfuscate their executable locations using common application folders that are not typically scrutinized. Additionally, it assumes that system monitoring of registry changes is enabled and functioning correctly. There is also a possibility of false positives from legitimate software applications that use unusual folders for updates or configurations.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
