# sigma/rules/windows/registry/registry_set/registry_set_asep_reg_keys_modification_system_scripts.yml
title: System Scripts Autorun Keys Modification
id: e7a2fd40-3ae1-4a85-bf80-15cf624fb1b1
related:
- id: 17f878b8-9968-4578-b814-c4217fc5768c
  type: obsolete
status: test
description: Detects modification of autostart extensibility point (ASEP) in registry.
references:
- https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1547.001/T1547.001.md
- https://learn.microsoft.com/en-us/sysinternals/downloads/autoruns
- https://gist.github.com/GlebSukhodolskiy/0fc5fa5f482903064b448890db1eaf9d     # a list with registry keys
author: Victor Sergeev, Daniil Yugoslavskiy, Gleb Sukhodolskiy, Timur Zinniatullin,
  oscd.community, Tim Shelton, frack113 (split)
date: 2019-10-25
modified: 2023-08-17
tags:
- attack.persistence
- attack.t1547.001
logsource:
  category: registry_set
  product: windows
detection:
  scripts_base:
    TargetObject|contains: '\Software\Policies\Microsoft\Windows\System\Scripts'
  scripts:
    TargetObject|contains:
    - '\Startup'
    - '\Shutdown'
    - '\Logon'
    - '\Logoff'
  filter:
    Details: '(Empty)'
  condition: scripts_base and scripts and not filter
fields:
- SecurityID
- ObjectName
- OldValueType
- NewValueType
falsepositives:
- Legitimate software automatically (mostly, during installation) sets up autorun
  keys for legitimate reason
- Legitimate administrator sets up autorun keys for legitimate reason
level: medium
notes: |
  ### Technical Context
  The Sigma rule “System Scripts Autorun Keys Modification” is designed to detect unauthorized modifications to the Autostart Extensibility Points (ASEPs) in the Windows registry. ASEPs are registry locations that allow executables and scripts to run automatically during system events such as startup or shutdown. This rule specifically looks for changes to keys under `\Software\Policies\Microsoft\Windows\System\Scripts`, particularly those related to autorun scripts for startup, shutdown, logon, and logoff processes. The collection of data sources essential for this detection includes registry modification logs, specifically monitoring changes to the TargetObject and identifying old and new values associated with these registry keys. 
  The rule aligns with the MITRE ATT&CK tactic of Persistence, particularly using technique T1547.001, which involves persistence mechanisms that allow attackers to maintain their foothold in the environment. By capturing modifications to these registry keys, organizations can uncover potentially malicious attempts to execute scripts without legitimate user consent.
  ### Investigation Steps
  - **Check EDR Endpoint Activity**: Review the endpoint's activity logs in the EDR for the processes that initiated the registry change, noting any associated command-line arguments or scripts involved.
  - **Review Registry Changes**: Use the built-in Windows Registry Editor or a tool like Sysinternals Autoruns to check the current and previous entries of the relevant autorun keys, comparing the new values against known legitimate configurations.
  - **Correlate User Actions**: Investigate the user account associated with the modification in the Windows Security logs to determine if the change was made by an authorized user or software agent.
  - **Analyze Network Traffic**: If external malicious command and control server communication is suspected, review network connection logs and proxy logs to identify any suspicious outbound traffic initiated from the affected endpoint around the time of the registry change.
  ### Prioritization
  This alert has a medium severity level because while changes to autorun keys can indicate malicious persistence attempts, they can also result from legitimate software installations or administrator actions. However, given that unauthorized modifications often play a key role in malware execution, further investigation is warranted.
  ### Blind Spots and Assumptions
  This rule may fail to fire in scenarios where modifications to the autorun keys do not generate registry events due to silent installations by certain legitimate applications or when malicious actors utilize stealth techniques to alter autorun scripts without direct registry modification. Additionally, if an organization heavily relies on legitimate administrative tasks that change these keys, there could be an increased chance of false positives. It is assumed that the monitoring systems are adequately configured to capture all relevant registry write events, and analysts should be aware of the valid software that could trigger these alerts.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
