# sigma/rules/windows/registry/registry_set/registry_set_disable_windows_event_log_access.yml
title: Windows Event Log Access Tampering Via Registry
id: ba226dcf-d390-4642-b9af-b534872f1156
status: experimental
description: |
  Detects changes to the Windows EventLog channel permission values. It focuses on changes to the Security Descriptor Definition Language (SDDL) string, as modifications to these values can restrict access to specific users or groups, potentially aiding in defense evasion by controlling who can view or modify a event log channel. Upon execution, the user shouldn't be able to access the event log channel via the event viewer or via utilities such as "Get-EventLog" or "wevtutil".
references:
- https://www.atomicredteam.io/atomic-red-team/atomics/T1562.002#atomic-test-8---modify-event-log-channel-access-permissions-via-registry---powershell
- https://www.youtube.com/watch?v=uSYvHUVU8xY
- https://learn.microsoft.com/en-us/windows/win32/secauthz/security-descriptor-definition-language
author: X__Junior
date: 2025-01-16
modified: 2025-02-05
tags:
- attack.t1547.001
- attack.t1112
logsource:
  category: registry_set
  product: windows
detection:
    # O:SYG:SYD:(D;;0x1;;;WD)
    # O:BAG:SYD:(A;;0x2;;;S-1-15-2-1)(D;;0x1;;;WD)
  selection_key_1:
    TargetObject|contains: '\SYSTEM\CurrentControlSet\Services\EventLog\'
    TargetObject|endswith: '\CustomSD'
  selection_key_2:
    TargetObject|contains:
    - '\Policies\Microsoft\Windows\EventLog\'
    - '\Microsoft\Windows\CurrentVersion\WINEVT\Channels'
    TargetObject|endswith: '\ChannelAccess'
  selection_details:
  - Details|contains: 'D:(D;'
  - Details|contains|all:
    - 'D:('
    - ')(D;'
  condition: 1 of selection_key_* and selection_details
falsepositives:
- Administrative activity, still unlikely
level: high
notes: |
  ### Technical Context
  This Sigma rule focuses on detecting unauthorized changes to Windows Event Log channel permissions via modifications in the registry. It specifically examines alterations to the Security Descriptor Definition Language (SDDL) string within the registry entries associated with the Windows Event Log. By monitoring keys related to event logging permissions (e.g., `\SYSTEM\CurrentControlSet\Services\EventLog\`, `\Policies\Microsoft\Windows\EventLog\`, and `\Microsoft\Windows\CurrentVersion\WINEVT\Channels`), the rule identifies attempts to restrict access to event logs, which may facilitate evasion techniques utilized by attackers. Such modifications can prevent legitimate users from viewing or accessing event logs through tools like "Get-EventLog" or "wevtutil", thereby hiding malicious activities. This behavior is mapped to the MITRE ATT&CK techniques for "Modify Event Log Channel Access Permissions" (T1562.002) and "Impair Defenses" (T1112).
  ### Investigation Steps
  - **Check Registry Alterations:** Use EDR solutions to query recent registry changes under the specified keys to identify any unauthorized modifications to event log permissions.
  - **Review Event Logs:** Inspect Windows Security and System event logs for unusual activities or events preceding or following the detection of the registry modifications to correlate with administrative actions.
  - **Analyze User Activity:** Cross-reference user activity via logins and actions against the altered permissions, determining if unauthorized users had access before the changes.
  - **Engage Incident Response:** If malicious changes are confirmed, escalate the situation through the incident response team to contain and remediate potential security breaches.
  ### Prioritization
  Given the potential for such changes to conceal malicious activity, the severity of this alert is classified as high. Unauthorized alterations to event log permissions can significantly hinder incident detection and response efforts, making it critical for immediate investigation.
  ### Blind Spots and Assumptions
  This rule may not fire in scenarios where legitimate administrative actions modify the Event Log permissions; therefore, it is crucial to assess context around the event. Additionally, an adversary could employ stealthy methods to change event log permissions without detection, especially if they have elevated privileges. It's essential for detection engineers to regularly update and validate the effectiveness of this rule to account for new evasion techniques and shifts in administrative practices.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
