# sigma/rules/windows/registry/registry_set/registry_set_asep_reg_keys_modification_winsock2.yml
title: WinSock2 Autorun Keys Modification
id: d6c2ce7e-afb5-4337-9ca4-4b5254ed0565
related:
- id: 17f878b8-9968-4578-b814-c4217fc5768c
  type: derived
status: test
description: Detects modification of autostart extensibility point (ASEP) in registry.
references:
- https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1547.001/T1547.001.md
- https://learn.microsoft.com/en-us/sysinternals/downloads/autoruns
- https://gist.github.com/GlebSukhodolskiy/0fc5fa5f482903064b448890db1eaf9d     # a list with registry keys
author: Victor Sergeev, Daniil Yugoslavskiy, Gleb Sukhodolskiy, Timur Zinniatullin,
  oscd.community, Tim Shelton, frack113 (split)
date: 2019-10-25
modified: 2023-08-17
tags:
- attack.persistence
- attack.t1547.001
logsource:
  category: registry_set
  product: windows
detection:
  winsock_parameters_base:
    TargetObject|contains: '\System\CurrentControlSet\Services\WinSock2\Parameters'
  winsock_parameters:
    TargetObject|contains:
    - '\Protocol_Catalog9\Catalog_Entries'
    - '\NameSpace_Catalog5\Catalog_Entries'
  filter:
  - Details: '(Empty)'
  - Image: 'C:\Windows\System32\MsiExec.exe'
  - Image: 'C:\Windows\syswow64\MsiExec.exe'
  condition: winsock_parameters_base and winsock_parameters and not filter
fields:
- SecurityID
- ObjectName
- OldValueType
- NewValueType
falsepositives:
- Legitimate software automatically (mostly, during installation) sets up autorun
  keys for legitimate reason
- Legitimate administrator sets up autorun keys for legitimate reason
level: medium
notes: |
  ### Technical Context
  The Sigma rule titled "WinSock2 Autorun Keys Modification" is designed to detect unauthorized modifications to Windows registry autorun keys associated with the WinSock2 service. This rule focuses on changes to specific registry paths that control how Windows loads network protocols and their related parameters. The key targets are the following registry locations: `\System\CurrentControlSet\Services\WinSock2\Parameters`, `\Protocol_Catalog9\Catalog_Entries`, and `\NameSpace_Catalog5\Catalog_Entries`. By monitoring modifications to these registry entries, the rule aims to identify potential persistence mechanisms employed by threat actors, particularly those related to the MITRE ATT&CK technique T1547.001 (Registry Run Keys / Startup Folder).
  The detection leverages registry set logs from Windows, specifically focusing on registry changes that involve the aforementioned keys. The rule applies filters to reduce false positives, excluding details like empty values or modifications triggered by legitimate processes, such as `MsiExec.exe`, which is used during software installations. As unauthorized changes to autorun keys can indicate an attempt by malware to persist on the system, this rule is valuable for SOC analysts and incident responders tasked with identifying such threats.
  ### Investigation Steps
  - **Review Registry Changes:** Use tools such as Sysinternals' Autoruns or PowerShell cmdlets to examine the specific registry paths for recent changes, focusing on the values set for the keys involved.
    
  - **Cross-Reference with Process Logs:** Utilize EDR solutions to check for process creation logs corresponding to the times when changes were made, specifically looking for processes that may have interacted with the `MsiExec.exe` images.
  - **Analyze User Context:** Investigate the Security ID (SID) from the registry change events to determine what account executed the modifications and whether it aligns with expected administrative activities.
  - **Check for Related Communications:** Monitor network connections and DNS queries associated with the affected system to identify any unusual outbound traffic that may indicate command and control (C2) communication or data exfiltration.
  ### Prioritization
  The alert generated by this rule is considered medium severity due to the potential implications of unauthorized registry modifications. A detected change could indicate an attempt by malicious actors to establish persistence on the compromised system, thereby warranting prompt investigation.
  ### Blind Spots and Assumptions
  This rule may not fire under certain conditions, such as legitimate software modifying autorun keys during installation processes, which are frequent and expected in a typical enterprise environment. Additionally, if an adversary uses techniques to alter their process image signature to appear legitimate, or if they modify the registry in a manner that does not trigger the detection parameters, the alert may not be generated. Analysts should be aware that not all modifications to these paths are malicious, and thorough investigation is necessary to validate the context of each detection.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
