# sigma/rules/windows/registry/registry_set/registry_set_timeproviders_dllname.yml
title: New TimeProviders Registered With Uncommon DLL Name
id: e88a6ddc-74f7-463b-9b26-f69fc0d2ce85
status: test
description: |
  Detects processes setting a new DLL in DllName in under HKEY_LOCAL_MACHINE\SYSTEM\CurrentControlSet\Services\W32Time\TimeProvider.
  Adversaries may abuse time providers to execute DLLs when the system boots.
  The Windows Time service (W32Time) enables time synchronization across and within domains.
references:
- https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1547.003/T1547.003.md
author: frack113
date: 2022-06-19
modified: 2024-03-26
tags:
- attack.persistence
- attack.privilege-escalation
- attack.t1547.003
logsource:
  category: registry_set
  product: windows
detection:
  selection:
    TargetObject|contains: '\Services\W32Time\TimeProviders'
    TargetObject|endswith: '\DllName'
  filter_main_w32time:
    Details:
    - '%SystemRoot%\System32\vmictimeprovider.dll'
    - '%systemroot%\system32\w32time.dll'
    - 'C:\Windows\SYSTEM32\w32time.DLL'
  condition: selection and not 1 of filter_main_*
falsepositives:
- Unknown
level: high
notes: |
  ### Technical Context
  The detection rule identified as "New TimeProviders Registered With Uncommon DLL Name" is designed to monitor changes to the Windows registry, specifically focusing on the `DllName` entry under the `HKEY_LOCAL_MACHINE\SYSTEM\CurrentControlSet\Services\W32Time\TimeProvider` path. This registry key is associated with the Windows Time service (W32Time), which is responsible for synchronizing time across systems in a domain. The detection rule looks for new DLLs registered at this path that do not match known safe DLLs, such as `vmictimeprovider.dll` and `w32time.dll`. The rule leverages Windows registry set logs to identify potential malicious activity where adversaries could exploit the time provider mechanism to execute unauthorized code during system boot. 
  This activity is linked to the MITRE ATT&CK technique T1547.003, which falls under the tactics of persistence and privilege escalation. By monitoring this specific registry change, the SOC can detect attempts to manipulate system time settings for greater control or evasion by threat actors.
  ### Investigation Steps
  - **Check EDR Alerts:** Review events related to processes utilizing the Windows Time service in your EDR solution to determine if any associated processes were malicious or unexpected.
  - **Examine Registry Changes:** Inspect the registry logs for recent modifications under the time provider path to identify any newly registered DLLs. Look for any entries that deviate from the norm or are not typically associated with the Windows Time service.
  - **Analyze Process Creation Logs:** Investigate process creation logs to identify which processes were spawned in conjunction with the registry changes. Focus on verifying the legitimacy and context of these processes.
  - **Review System and Application Logs:** Utilize Windows Event Logs to look for any abnormal behavior or system events around the time of the registry change. Pay particular attention to the System and Application logs for error messages or anomalies that could indicate misuse.
  ### Prioritization
  Given the high severity level assigned to this alert, it indicates a significant risk that adversaries may be attempting to establish persistence on the system. Prompt investigation is warranted to mitigate any potential escalation or breach of sensitive systems.
  ### Blind Spots and Assumptions
  This rule assumes that all registry changes are logged and that the telemetry is correctly configured to capture Windows registry modifications. However, it may fail to trigger if the adversary uses legitimate processes to manipulate the registry without raising alarms or if they utilize known safe DLLs that bypass the detection criteria. Additionally, if a system is poorly monitored or if security logging is not appropriately enabled, this rule may not fire, resulting in missed detections.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
