# sigma/rules/windows/registry/registry_set/registry_set_asep_reg_keys_modification_wow6432node_classes.yml
title: Wow6432Node Classes Autorun Keys Modification
id: 18f2065c-d36c-464a-a748-bcf909acb2e3
related:
- id: 17f878b8-9968-4578-b814-c4217fc5768c
  type: obsolete
status: test
description: Detects modification of autostart extensibility point (ASEP) in registry.
references:
- https://github.com/redcanaryco/atomic-red-team/blob/f339e7da7d05f6057fdfcdd3742bfcf365fee2a9/atomics/T1547.001/T1547.001.md
- https://learn.microsoft.com/en-us/sysinternals/downloads/autoruns
- https://gist.github.com/GlebSukhodolskiy/0fc5fa5f482903064b448890db1eaf9d     # a list with registry keys
author: Victor Sergeev, Daniil Yugoslavskiy, Gleb Sukhodolskiy, Timur Zinniatullin,
  oscd.community, Tim Shelton, frack113 (split)
date: 2019-10-25
modified: 2023-08-17
tags:
- attack.persistence
- attack.t1547.001
logsource:
  category: registry_set
  product: windows
detection:
  wow_classes_base:
    TargetObject|contains: '\Software\Wow6432Node\Classes'
  wow_classes:
    TargetObject|contains:
    - '\Folder\ShellEx\ExtShellFolderViews'
    - '\Folder\ShellEx\DragDropHandlers'
    - '\Folder\ShellEx\ColumnHandlers'
    - '\Directory\Shellex\DragDropHandlers'
    - '\Directory\Shellex\CopyHookHandlers'
    - '\CLSID\{AC757296-3522-4E11-9862-C17BE5A1767E}\Instance'
    - '\CLSID\{ABE3B9A4-257D-4B97-BD1A-294AF496222E}\Instance'
    - '\CLSID\{7ED96837-96F0-4812-B211-F13C24117ED3}\Instance'
    - '\CLSID\{083863F1-70DE-11d0-BD40-00A0C911CE86}\Instance'
    - '\AllFileSystemObjects\ShellEx\DragDropHandlers'
    - '\ShellEx\PropertySheetHandlers'
    - '\ShellEx\ContextMenuHandlers'
  filter:
    Details: '(Empty)'
  condition: wow_classes_base and wow_classes and not filter
fields:
- SecurityID
- ObjectName
- OldValueType
- NewValueType
falsepositives:
- Legitimate software automatically (mostly, during installation) sets up autorun
  keys for legitimate reason
- Legitimate administrator sets up autorun keys for legitimate reason
level: medium
notes: |
  n
  ### Technical Context
  The "Wow6432Node Classes Autorun Keys Modification" detection rule monitors for modifications within specific registry keys associated with the Windows Registry's Wow6432Node, which is utilized for 32-bit applications on 64-bit Windows systems. This rule primarily targets changes to authorun extensibility points (ASEP) that are commonly exploited by malicious actors to ensure persistence on compromised systems. It inspects registry changes to various subkeys related to shell extensions, such as `ShellEx` handlers for folders and directories, which are critical for enforcing autorun behaviors. The primary log source involved in this detection is the Windows registry change logs, specifically focused on the registry categories where autorun configurations are set.
  This rule aligns with the MITRE ATT&CK framework under the tactic of Persistence (TA0003) and specifically targets technique T1547.001, which relates to the modification of registry keys to establish persistence through autorun features. By capturing these modifications, the rule assists incident responders in identifying potential malicious activities aimed at maintaining unauthorized access to systems.
  ### Investigation Steps
  - **Review Registry Changes:** Use the EDR to review recent registry modifications, particularly around `HKEY_LOCAL_MACHINE\Software\Wow6432Node\Classes`, to identify the instigating process and timestamp of the changes.
  - **Correlate Processes:** Investigate any suspicious processes that executed around the time of the registry modification by examining process creation logs in Sysmon or Windows Event Logs.
  - **Query Autorun Entries:** Utilize Sysinternals Autoruns tool to collect detailed information about startup entries and correlate them against the suspicious changes detected, focusing on entries modified around the time of the alert.
  - **Check for Legitimate Reasons:** Identify if the registry changes were performed by authorized administrators or legitimate software installations by analyzing user activity logs and application installation events.
  ### Prioritization
  This alert is categorized as medium severity because while registry modifications can occasionally be legitimate, unauthorized changes in ASEP registry keys are often indicative of a persistence tactic used by attackers, warranting immediate investigative follow-up.
  ### Blind Spots and Assumptions
  This rule may not fire in scenarios where malicious actors utilize signed or trusted applications to modify autorun keys without detection. Additionally, changes made in the registry by legitimate software during installation can create false positives, particularly in environments where software installations occur frequently. There may also be variations in registry access permissions or policies across different environments that could lead to undetected modifications. Therefore, it is prudent to ensure proper endpoint auditing and to understand the software baseline expected within the organization to effectively narrow down the investigation scope.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
