# sigma/rules/windows/builtin/security/win_security_susp_scheduled_task_update.yml
title: Suspicious Scheduled Task Update
id: 614cf376-6651-47c4-9dcc-6b9527f749f4
related:
- id: 1c0e41cd-21bb-4433-9acc-4a2cd6367b9b     # ProcCreation schtasks change
  type: similar
status: test
description: Detects update to a scheduled task event that contain suspicious keywords.
references:
- https://learn.microsoft.com/en-us/previous-versions/windows/it-pro/windows-10/security/threat-protection/auditing/event-4698
author: Nasreddine Bencherchali (Nextron Systems)
date: 2022-12-05
tags:
- attack.execution
- attack.privilege-escalation
- attack.persistence
- attack.t1053.005
logsource:
  product: windows
  service: security
  definition: 'The Advanced Audit Policy setting Object Access > Audit Other Object
    Access Events has to be configured to allow this detection. We also recommend
    extracting the Command field from the embedded XML in the event data.'
detection:
  selection_eid:
    EventID: 4702
  selection_paths:
    TaskContentNew|contains:
    - '\AppData\Local\Temp\'
    - '\AppData\Roaming\'
    - '\Users\Public\'
    - '\WINDOWS\Temp\'
    - 'C:\Temp\'
    - '\Desktop\'
    - '\Downloads\'
    - '\Temporary Internet'
    - 'C:\ProgramData\'
    - 'C:\Perflogs\'
  selection_commands:
    TaskContentNew|contains:
    - 'regsvr32'
    - 'rundll32'
    - 'cmd.exe</Command>'
    - 'cmd</Command>'
    - '<Arguments>/c '
    - '<Arguments>/k '
    - '<Arguments>/r '
    - 'powershell'
    - 'pwsh'
    - 'mshta'
    - 'wscript'
    - 'cscript'
    - 'certutil'
    - 'bitsadmin'
    - 'bash.exe'
    - 'bash '
    - 'scrcons'
    - 'wmic '
    - 'wmic.exe'
    - 'forfiles'
    - 'scriptrunner'
    - 'hh.exe'
  condition: all of selection_*
falsepositives:
- Unknown
level: high
notes: |
  n
  ### Technical Context
  The "Suspicious Scheduled Task Update" detection rule monitors Windows Security event logs for updates to scheduled tasks that are associated with potentially malicious activity. Specifically, it looks for Event ID 4702, which signals the modification of scheduled tasks, and scans the task's content for suspicious keywords or paths typically associated with exploitation attempts. This includes examining command-line arguments for known triggering executables like `powershell`, `cmd.exe`, or `regsvr32`, as well as directories like `Temp`, `Roaming`, and `Public`. These enhancements are important for typical strategies adversaries may employ in persistence and privilege escalation, in line with the MITRE ATT&CK technique T1053.005 — Scheduled Task/Job.
  ### Investigation Steps
  - **Verify Task Modifications:** Use EDR tools to examine the configured tasks for the event in question. Compare the task settings with baseline configurations or known-good states to identify anomalies.
    
  - **Command Execution Review:** Analyze the command execution parameters from the task's content through the Security logs or Sysmon data. Look for execution of suspicious commands or scripts related to known exploits.
  - **Path Analysis:** Investigate the file system paths associated with the scheduled tasks. Check if the paths are legitimate locations or if they redirect to temporary or public directories that could be synonymous with malicious activities.
  - **Cross-Reference Alerts:** Utilize correlation features within the SIEM or EDR to check if this alert aligns with other suspicious activities, such as failed logins or known malware signatures, to assess the broader context of the incident.
  ### Prioritization
  Given the potential for privilege escalation and persistence that may accompany suspicious scheduled task modifications, alerts from this rule are classified as high severity. Quick identification and response are imperative to prevent further exploitation.
  ### Blind Spots and Assumptions
  This rule presumes that the environment has adequate logging configured for Event ID 4702 and may fail to trigger if the Advanced Audit Policy settings are not correctly implemented. Additionally, if an attacker has sufficient access, they may create or modify tasks without triggering alerts, especially if they utilize legitimate tools or scripts that do not meet the identified criteria. This means that continued vigilance and proactive configurations are necessary to maintain security posture.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
