# sigma/rules/windows/builtin/security/win_security_gpo_scheduledtasks.yml
title: Persistence and Execution at Scale via GPO Scheduled Task
id: a8f29a7b-b137-4446-80a0-b804272f3da2
status: test
description: Detect lateral movement using GPO scheduled task, usually used to deploy
  ransomware at scale
references:
- https://twitter.com/menasec1/status/1106899890377052160
- https://www.secureworks.com/blog/ransomware-as-a-distraction
- https://www.elastic.co/guide/en/security/7.17/prebuilt-rule-0-16-1-scheduled-task-execution-at-scale-via-gpo.html
author: Samir Bousseaden
date: 2019-04-03
modified: 2024-09-04
tags:
- attack.persistence
- attack.lateral-movement
- attack.t1053.005
logsource:
  product: windows
  service: security
  definition: 'The advanced audit policy setting "Object Access > Audit Detailed File
    Share" must be configured for Success/Failure'
detection:
  selection_5136:
    EventID: 5136
    AttributeLDAPDisplayName:
    - 'gPCMachineExtensionNames'
    - 'gPCUserExtensionNames'
    AttributeValue|contains:
    - 'CAB54552-DEEA-4691-817E-ED4A4D1AFC72'
    - 'AADCED64-746C-4633-A97C-D61349046527'
  selection_5145:
    EventID: 5145
    ShareName|endswith: '\SYSVOL'     # looking for the string \\*\SYSVOL
    RelativeTargetName|endswith: 'ScheduledTasks.xml'
    AccessList|contains:
    - 'WriteData'
    - '%%4417'
  condition: 1 of selection_*
falsepositives:
- If the source IP is not localhost then it's super suspicious, better to monitor
  both local and remote changes to GPO scheduled tasks.
level: high
notes: |
  ### Technical Context
  This Sigma rule focuses on detecting lateral movement techniques that utilize Group Policy Object (GPO) scheduled tasks, particularly in the context of attacks like ransomware deployment at scale. Attacks leveraging GPO scheduled tasks can change task settings or introduce new tasks via LDAP modifications, especially targeting the `gPCMachineExtensionNames` and `gPCUserExtensionNames` attributes of GPOs. The rule monitors events indicated by Event IDs 5136 and 5145 in Windows Security logs, specifically looking for unauthorized changes to properties related to scheduled tasks and anticipated access rights to shared resources. These changes can escalate existing privileges or allow for the concise execution of malicious scripts across the network, fulfilling the MITRE ATT&CK techniques for persistence (T1053.005) and lateral movement.
  ### Investigation Steps
  - **Check Windows Security Logs**: Look for Event ID 5136 and 5145 that indicate unauthorized changes to GPO scheduled tasks, focusing on details of the modifications under scrutiny. Validate whether the changes were performed from localhost or are unusually initiated from remote systems.
  - **Analyze Active Directory Changes**: Use tools like PowerShell or Active Directory management consoles to investigate recent modifications to GPOs, focusing on the attributes mentioned in the rule. Confirm that the changes correspond with legitimate administrative activities.
  - **Inspect Network Activity**: Utilize EDR and NDR tools to examine network traffic from the originating device to identify any suspicious communication attempts to other endpoints, especially those targeting resources that might execute tasks under altered GPOs.
  - **Monitor for Related Alerts**: Cross-reference alerts from EDRs, AV, and other security solutions to check for known indicators of compromise or abnormal behaviors following the changes detected by the rule. Ensure no secondary malicious payloads were deployed concurrently.
  ### Prioritization
  The alert generated by this detection rule is considered high severity due to the critical nature of GPO modifications, which can signify a deep compromise of the network and allow for extensive lateral movement. Immediate investigation is warranted to mitigate potential widespread impact.
  ### Blind Spots and Assumptions
  This rule may not trigger if attackers use legitimate accounts to make GPO changes, or if detailed auditing is not enabled for file shares targeting the SYSVOL directory. Furthermore, it assumes that all necessary event IDs are consistently logged and accessible in the environment. Engineers should consider potential evasion tactics such as using alternate accounts or tools that can modify GPOs without triggering the specified events or those that filter out audit records to hide malicious actions.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
