# sigma/rules/windows/builtin/security/win_security_susp_scheduled_task_delete_or_disable.yml
title: Important Scheduled Task Deleted/Disabled
id: 7595ba94-cf3b-4471-aa03-4f6baa9e5fad
related:
- id: dbc1f800-0fe0-4bc0-9c66-292c2abe3f78     # ProcCreation schtasks delete
  type: similar
- id: 9ac94dc8-9042-493c-ba45-3b5e7c86b980     # ProcCreation schtasks disable
  type: similar
- id: 9e3cb244-bdb8-4632-8c90-6079c8f4f16d     # TaskScheduler EventLog
  type: similar
status: test
description: Detects when adversaries stop services or processes by deleting or disabling
  their respective scheduled tasks in order to conduct data destructive activities
references:
- https://learn.microsoft.com/en-us/previous-versions/windows/it-pro/windows-10/security/threat-protection/auditing/event-4699
- https://learn.microsoft.com/en-us/previous-versions/windows/it-pro/windows-10/security/threat-protection/auditing/event-4701
author: Nasreddine Bencherchali (Nextron Systems)
date: 2022-12-05
modified: 2023-03-13
tags:
- attack.execution
- attack.privilege-escalation
- attack.persistence
- attack.t1053.005
logsource:
  product: windows
  service: security
  definition: 'The Advanced Audit Policy setting Object Access > Audit Other Object
    Access Events has to be configured to allow this detection. We also recommend
    extracting the Command field from the embedded XML in the event data.'
detection:
  selection:
    EventID:
    - 4699         # Task Deleted Event
    - 4701         # Task Disabled Event
    TaskName|contains:
            # Add more important tasks
    - '\Windows\SystemRestore\SR'
    - '\Windows\Windows Defender\'
    - '\Windows\BitLocker'
    - '\Windows\WindowsBackup\'
    - '\Windows\WindowsUpdate\'
    - '\Windows\UpdateOrchestrator\Schedule'
    - '\Windows\ExploitGuard'
  filter_sys_username:
    EventID: 4699
    SubjectUserName|endswith: '$'      # False positives during upgrades of Defender, where its tasks get removed and added
    TaskName|contains: '\Windows\Windows Defender\'
  condition: selection and not 1 of filter_*
falsepositives:
- Unknown
level: high
notes: |
  ### Technical Context
  This detection rule identifies when critical Windows scheduled tasks are either deleted or disabled, which may indicate malicious activity intended to halt security measures or disrupt system operations. Specifically, the rule monitors Windows Security event logs for EventIDs 4699 (Task Deleted) and 4701 (Task Disabled). It targets a predefined list of vital scheduled tasks, such as those related to system restore, Windows Defender, BitLocker, and Windows Update, which are essential for maintaining system integrity and security. When an adversary seeks to undermine system defenses or execute data-destructive activities, the deletion or disabling of these scheduled tasks can serve as a flag. This rule aligns with the MITRE ATT&CK technique T1053.005 (Scheduled Task/Job: Delete Scheduled Task).
  ### Investigation Steps
  - **Review Event Logs:** Check Windows Security event logs for EventIDs 4699 and 4701, focusing on the timeframe around when the alert was triggered to identify the specific tasks affected.
  - **Correlate with Process Creation:** Utilize EDR tools to correlate the deletion or disabling events with any suspicious process creation (e.g., processes invoking `schtasks.exe`) that occurred in proximity to these events to assess whether a malicious actor orchestrated these changes.
  - **Check User Context:** Investigate the user context under which the tasks were deleted or disabled, specifically looking for accounts that should not have performed these actions (such as non-privileged users).
  - **Examine Threat Intelligence:** Consult threat intelligence sources to see if the identified tasks are commonly targeted by known threat actors, which could provide additional context for the incident response.
  ### Prioritization
  This alert is classified as high severity due to the critical nature of the scheduled tasks involved and the potential for serious ramifications if an adversary successfully disables essential system processes, which can lead to broader system compromise.
  ### Blind Spots and Assumptions
  This rule assumes that all relevant scheduled tasks are logged and monitored correctly, and may not trigger if a legitimate task is renamed or moved without deletion or disabling. If attackers utilize privileged accounts or employ tactics to mimic legitimate administrative activity, the rule might fail to catch their actions. Additionally, if the necessary auditing policies are not correctly configured, scheduled task modifications might not be logged.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
