# sigma/rules/windows/builtin/security/win_security_susp_scheduled_task_creation.yml
title: Suspicious Scheduled Task Creation
id: 3a734d25-df5c-4b99-8034-af1ddb5883a4
status: test
description: Detects suspicious scheduled task creation events. Based on attributes
  such as paths, commands line flags, etc.
references:
- https://learn.microsoft.com/en-us/previous-versions/windows/it-pro/windows-10/security/threat-protection/auditing/event-4698
author: Nasreddine Bencherchali (Nextron Systems)
date: 2022-12-05
modified: 2022-12-07
tags:
- attack.execution
- attack.privilege-escalation
- attack.persistence
- attack.t1053.005
logsource:
  product: windows
  service: security
  definition: 'The Advanced Audit Policy setting Object Access > Audit Other Object
    Access Events has to be configured to allow this detection. We also recommend
    extracting the Command field from the embedded XML in the event data.'
detection:
  selection_eid:
    EventID: 4698
  selection_paths:
    TaskContent|contains:
    - '\AppData\Local\Temp\'
    - '\AppData\Roaming\'
    - '\Users\Public\'
    - '\WINDOWS\Temp\'
    - 'C:\Temp\'
    - '\Desktop\'
    - '\Downloads\'
    - '\Temporary Internet'
    - 'C:\ProgramData\'
    - 'C:\Perflogs\'
  selection_commands:
    TaskContent|contains:
    - 'regsvr32'
    - 'rundll32'
    - 'cmd.exe</Command>'
    - 'cmd</Command>'
    - '<Arguments>/c '
    - '<Arguments>/k '
    - '<Arguments>/r '
    - 'powershell'
    - 'pwsh'
    - 'mshta'
    - 'wscript'
    - 'cscript'
    - 'certutil'
    - 'bitsadmin'
    - 'bash.exe'
    - 'bash '
    - 'scrcons'
    - 'wmic '
    - 'wmic.exe'
    - 'forfiles'
    - 'scriptrunner'
    - 'hh.exe'
  condition: all of selection_*
falsepositives:
- Unknown
level: high
notes: |
  ### Technical Context
  This Sigma rule detects suspicious scheduled task creation events in Windows environments by monitoring the Security Event Log for Event ID 4698. The rule targets specific attributes within the task content, such as file paths and command-line flags often associated with malicious activity. It particularly examines paths frequently used by attackers (e.g., temporary folders and the AppData directory) and commands that can invoke potentially harmful executables (e.g., `regsvr32`, `powershell`, and `cmd.exe`). The participation of the Advanced Audit Policy setting for Object Access is essential for accurate detection, allowing analysts to glean detailed information from the event data. This detection aligns with the MITRE ATT&CK tactics of Execution (T1053.005), Privilege Escalation, and Persistence, making it crucial for identifying activities that may indicate an attacker's persistence mechanisms.
  ### Investigation Steps
  - **Review Event Logs:** Check the Security Event Log for Event ID 4698 to identify details of the suspicious scheduled task, including its creation time and user context.
  - **Analyze Task Configuration:** Use tools like PowerShell or the Task Scheduler GUI to inspect the properties of the scheduled task for any abnormal execution parameters or paths that may indicate malicious intent.
  - **Correlate with EDR:** Utilize EDR tools to examine related process activity and ensure the task's creation wasn't preceded or followed by other suspicious actions often associated with malware.
  - **Check File System Activity:** Review file system logs for access or modifications to the suspicious directories mentioned in the task content to identify any other malicious artifacts or behaviors.
  ### Prioritization
  Given that this alert highlights potentially malicious scheduled task creations, it is classified as high severity due to the implications of persistence and privilege escalation, which could lead to further compromise if not addressed promptly.
  ### Blind Spots and Assumptions
  There are potential blind spots if the Advanced Audit Policy for Object Access is not properly configured, leading to missed scheduled task events. Additionally, not all suspicious tasks may have clearly malicious intentions, resulting in false positives which analysts must investigate further. Adversaries may defeat this detection tactic by using legitimate tools in permitted paths or employing obfuscation techniques in their scheduled tasks.
  > **Disclaimer:** This investigation guide was created using generative AI technology and has not been reviewed for its accuracy and relevance. While every effort has been made to ensure its quality, we recommend validating the content and adapting it to suit specific environment and operational needs. Please communicate any changes to the detection engineering team.
